# desktop/system/systemd/user/awesome-lock.target
[Unit]
Description=awesome wm session lock

BindsTo=i3lock.service
After=i3lock.service

