# desktop/system/systemd/user/console.target
[Unit]
Description=console configuration
AllowIsolate=true

Conflicts=graphical-session.target
Before=graphical-session.target
