# desktop/system/systemd/user/sway-session.target
[Unit]
Description=sway wm session
AllowIsolate=true

Wants=dunst.service xrdb.service
