# desktop/system/systemd/user/sway-lock.target
[Unit]
Description=sway wm session lock

BindsTo=swaylock.service
After=swaylock.service

