# desktop/system/systemd/user/xserver.target
[Unit]
Description=xwindow server

RefuseManualStart=true
StopWhenUnneeded=true

Before=graphical-session.target

Conflicts=console.target
After=console.target
