# user/hi.target.wants/initial-async.target
[Unit]
Description=initial synchronization

StopWhenUnneeded=true
AllowIsolate=true

BindsTo=async@default.service
After=async@default.service
