# user/hi.target.wants/initial-fetch.target
[Unit]
Description=initial mail fetch

StopWhenUnneeded=true
AllowIsolate=true

BindsTo=fetch.target
After=fetch.target
After=initial-sync.target
