# system/profile.d/fasd.sh

# load fasd hooks and aliases
eval "$(fasd --init auto)"
