# system/profile.d/xkb-layout.sh

# Configure keyboard layout for libxkbcommon
XKB_DEFAULT_MODEL="pc105"
XKB_DEFAULT_LAYOUT="es"
XKB_DEFAULT_VARIANT="cat"
