# user/sshmux@quark.service.d/openvpn.conf
[Unit]
Requires=openvpn@quark.service
After=openvpn@quark.service
