# user/sshmux@quark.service.d/wireguard.conf
[Unit]
Requires=wireguard@home.service
After=wireguard@home.service
