# user/gpgfwd@quark.service.d/openvpn.conf
[Unit]
BindsTo=openvpn@%i.service
After=openvpn@%i.service
