# user/sshmux@skynet.service.d/ec2.conf
[Unit]
Requires=ec2@skynet.service
After=ec2@skynet.service
