# user/sshmux@skynet.service.d/machine.conf
[Unit]
Requires=machine@skynet.service
After=machine@skynet.service
