# user/sshmux@skynet.service.d/openvpn.conf
[Unit]
Requires=openvpn@skynet.service
After=openvpn@skynet.service
