# systemd/user/chromium@.service
[Unit]
Description=chromium %i

Requisite=wm.target 
After=wm.target

[Service]
Type=simple
SyslogIdentifier=chromium
Slice=apps.slice

ExecStart=/usr/bin/chromium
StandardOutput=null


