# systemd/user/chromium@.service
[Unit]
Description=chromium %i

Requisite=wm.target 
After=wm.target

[Service]
Type=simple
Slice=apps.slice

ExecStart=/usr/bin/chromium --profile-directory=Default
StandardOutput=null


