# systemd/user/chromium@.service
[Unit]
Description=chromium %i

PartOf=graphical-session.target 
After=graphical-session.target

[Service]
Type=simple
Slice=apps.slice

ExecStart=/usr/bin/chromium --new-window --class=chromium-%I --user-data-dir=%h/var/chromium/%I
StandardOutput=null


