# systemd/user/clean-tmp.service
[Unit]
Description=clean tmp files

[Service]
Type=oneshot
Slice=workers.slice

SyslogIdentifier=clean-tmp
# NOTE: Do not use environment variables here. If they are empty, it can mess things up royally.
ExecStart=/usr/bin/find "%h/tmp/" -ctime +90 -mindepth 1 -exec rm -Rf {} +
ExecStart=/usr/bin/find "%h/opt/build/arch/" -ctime +10 -mindepth 1 -maxdepth 1 -exec rm -Rf {} +
ExecStart=/usr/bin/find "%h/opt/build/sage/" -ctime +10 -mindepth 1 -maxdepth 1 -exec rm -Rf {} +
