# systemd/user/clean-build.timer
[Unit]
Description=daily build cleanup

Requisite=housekeeping.target
After=housekeeping.target

[Timer]
OnCalendar=daily
RandomizedDelaySec=12h
Persistent=true
Unit=clean-build.service

[Install]
WantedBy=housekeeping.target