# systemd/user/default.target
[Unit]
Description=Default target
Requires=dbus.socket
AllowIsolate=true
