# systemd/user/default.target
[Unit]
Description=default user target
AllowIsolate=true
