# systemd/user/default.target
[Unit]
Description=default user target
AllowIsolate=true

Wants=user.target console.target
