# systemd/user/default.target
[Unit]
Description=user target
AllowIsolate=true

Wants=user.target console.target
