# systemd/user/deluge.service
[Unit]
Description=deluge client
After=graphical-session.target

PartOf=graphical-session.target 
After=graphical-session.target

Requires=deluged.service
After=deluged.service

[Service]
Type=simple

Slice=apps.slice

ExecStart=/usr/bin/deluge


