# systemd/user/desktop.target
[Unit]
Description=Desktop Environment
After=xorg@0.service
StopWhenUnneeded=true

Requires=dbus.socket
Requisite=xorg@0.service

AllowIsolate=true
