# systemd/user/device-impulse.target
[Unit]
Description=novation impulse device

BindsTo=sys-subsystem-sound-devices-impulse.device
After=sys-subsystem-sound-devices-impulse.device