# systemd/user/dpms-on.service
[Unit]
Description=dpms always powered

Requisite=graphical.target
After=graphical.target

[Service]
Type=oneshot
SyslogIdentifier=dpms-on
ExecStart=/usr/bin/xset dpms 0 0 0
ExecStart=/usr/bin/xset dpms force on