# systemd/user/dropbox.service
[Unit]
Description=dropbox daemon

After=xorg.target
Requisite=xorg.target

[Service]
Type=simple
ExecStart=/usr/bin/dropboxd
