# systemd/user/dropbox.service
[Unit]
Description=dropbox daemon

After=graphical.target
Requisite=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/dropboxd
