# systemd/user/dropbox.service
[Unit]
Description=dropbox daemon

After=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/dropboxd
