# systemd/user/emacsclient@.service
[Unit]
Description=emacs client %i

Requisite=wm.target 
After=wm.target

[Service]
Type=simple
Slice=apps.slice

ExecStart=/home/abdo/bin/emacs-client


