# systemd/user/fetch-contacts.target
[Unit]
Description=fetch contacts
StopWhenUnneeded=true

Requires=gcont-sync.service
After=gcont-sync.service
