# systemd/user/fetch.target
[Unit]
Description=fetch mail and news
Requires=dbus.socket
StopWhenUnneeded=true
AllowIsolate=true
Wants=offlineimap.service rss2maildir.service mutag.service

