# systemd/user/fetch-news.target
[Unit]
Description=fetch news
StopWhenUnneeded=true

Requires=mail-db.service
After=mail-db.service

Requires=mutag.service rss2maildir.service
After=mutag.service rss2maildir.service
