# systemd/user/final-synced-all.target
[Unit]
Description=final homedir sync

# refresh cached passphrase now
Requires=gpg-agent-auth.service

Conflicts=initial-synced-base.target initial-synced-lib.target initial-synced-mail.target
Conflicts=initial-synced-work.target initial-synced-rest.target

Requires=async-default-all.service
After=async-default-all.service

