# systemd/user/gdrive.service
[Unit]
Description=sync google drive

[Service]
Type=oneshot
ExecStart=/usr/bin/grive --path %h/share/gdrive

