# systemd/user/git-drop@.service
[Unit]
Description=drop unused annex files at /home/%u/%I

ConditionPathExists=/home/%u/%I/.git/annex

[Service]
Type=oneshot
Slice=workers.slice

ExecStart=/usr/bin/git -C "/home/%u/%I" annex unused --quiet
ExecStart=/usr/bin/git -C "/home/%u/%I" annex dropunused --quiet all
