# systemd/user/glib-pacrunner.service
[Unit]
Description=glib pacrunner

[Service]
Type=dbus
BusName=org.gtk.GLib.PACRunner
Slice=wm.slice

ExecStart=/usr/lib/glib-networking/glib-pacrunner

[Install]
Alias=dbus-org.gtk.GLib.PACRunner.service
