# systemd/user/gpg-agent.service
[Unit]
Description=gpg-agent daemon

ConditionPathExists=/home/%u/.gnupg

[Service]
Type=forking
SyslogIdentifier=gpg-agent
ExecStart=/usr/bin/gpg-agent --daemon --use-standard-socket

# give it some time
ExecStartPost=/usr/bin/sleep 1