# systemd/user/hodge-suspend.timer
[Unit]
Description=suspend hodge at night
ConditionHost=hodge

[Timer]
OnCalendar=20:00
Unit=byez.target

[Install]
WantedBy=timers.target