# systemd/user/homedir-gc.timer
[Unit]
Description=weekly garbage collection of homedir

Requisite=housekeeping.target

ConditionPathExists=/home/abdo/.mrconfig

[Timer]
OnCalendar=weekly
AccuracySec=12h
Persistent=true
Unit=homedir-gc.service

[Install]
WantedBy=housekeeping.target