# systemd/user/homedir-sync@.service
[Unit]
Description=sync homedir to remote %I

AssertPathExists=/home/%u/arch/sync/Makefile
AssertPathExists=/home/%u/priv

After=mount-priv.service

[Service]
Type=oneshot
Slice=workers.slice

WorkingDirectory=/home/%u/arch/sync
ExecStart=/usr/bin/make sync REMOTES=%I

