# systemd/user/homedir-gc@.service
[Unit]
Description=garbage collect directory /home/%u/%I

AssertPathExists=/home/%u/%I
AssertPathExists=/home/%u/arch/sync/Makefile

[Service]
Type=oneshot
Slice=workers.slice

WorkingDirectory=/home/%u/arch/sync
ExecStart=/usr/bin/make DIRS=%I gc
