# systemd/user/housekeeping.target
[Unit]
Description=houskeeping timers
PartOf=synced.target
