# systemd/user/housekeeping.target
[Unit]
Description=houskeeping timers

ConditionACPower=yes
Requisite=synced.target
