# systemd/user/i3lock.service
[Unit]
Description=i3lock screen locker

Requisite=graphical.target
After=graphical.target

# so that when i3lock ends lock.target goes down
PartOf=lock.target

[Service]
Type=simple
ExecStart=/usr/bin/i3lock --dpms -I 10 --nofork --color=000000

[Install]
WantedBy=lock.target
