# systemd/user/initial-synced-mail.target
[Unit]
Description=initial mail sync

Requires=async-default-mail.service
After=async-default-mail.service

