# systemd/user/journal-mailer.service
[Unit]
Description=monitor journal

Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
SyslogIdentifier=journal-mailer
ExecStart=/home/abdo/.config/systemd/scripts/journal-mailer.sh
