# systemd/user/journal-monitor.service
[Unit]
Description=Journal monitor

[Service]
Type=simple
SyslogIdentifier=journal-monitor
ExecStart=/home/abdo/.config/systemd/scripts/journal-monitor.py
