# systemd/user/kill-slice@.service
[Unit]
Description=kill %i slice

# NOTE: we need to manually kill services in slices, because transient units can't have
# Requires dependencies that would do that for us.

Before=wmclean.target

[Service]
Type=oneshot
Slice=workers.slice

ExecStart=/usr/bin/systemctl --user stop %i.slice
