# systemd/user/ladish.service
[Unit]
Description=ladish daemon

[Service]
Type=dbus
BusName=org.ladish
Slice=audio.slice

ExecStart=/usr/bin/ladishd

[Install]
Alias=dbus-org.ladish.service
