# systemd/user/mail.target
[Unit]
Description=fetch mail and launch mu4e

Wants=fetch.target 

BindsTo=mu4e.service
Before=mu4e.service

Requires=mail-db.service
After=mail-db.service

