# systemd/user/mpc-pause.service
[Unit]
Description=run mpc pause
ConditionPathExists=%t/mpd/pid

[Service]
Type=oneshot
ExecStart=/usr/bin/mpc pause
StandardOutput=null
