# systemd/user/mpc@.service
[Unit]
Description=Run mpc %i
After=mpd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mpc %i 
