# systemd/user/mpd.socket
[Unit]
Description=Music Player Daemon socket

[Socket]
ListenStream=%t/mpd/mpd.socket
SocketMode=0600
DirectoryMode=0700
Service=mpd.service

[Install]
WantedBy=synced.target
