# systemd/user/mpdscribble.service
[Unit]
Description=mpd scribbler

Requires=mpd.socket
After=mpd.service
Requires=mpd.service

[Service]
Type = simple
ExecStart = /usr/bin/mpdscribble -D
Restart = always

