# systemd/user/mpdscribble.service
[Unit]
Description=mpd scribbler

Requisite=synced.target
After=synced.target

Requires=mpd.socket mpd.service
After=mpd.socket mpd.service

Requires=audio.slice
After=audio.slice

[Service]
Type = simple
ExecStart = /usr/bin/mpdscribble -D

Slice=audio.slice