# systemd/user/mu-server.socket
[Unit]
Description=socket for mu server

[Socket]
ListenStream=%t/mu/mu.socket
SocketMode=0600

[Install]
WantedBy=sockets.target
