# systemd/user/music-synth.target
[Unit]
Description=play synthetizer

Wants=jack.service impulse.service
After=jack.service impulse.service

Wants=fluidsynth.service qjackctl-play.service
After=fluidsynth.service qjackctl-play.service
