# systemd/user/orgmode-db.service
[Unit]
Description=orgmode database

StopWhenUnneeded=true
ConditionPathExists=/home/%u/work/wiki

Requisite=synced.target
After=synced.target

[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=orgmode-db
ExecStart=/usr/bin/true
ExecStop=/usr/bin/systemctl --user start commit-orgmode.service


