# systemd/user/pacman-update.service
[Unit]
Description=update packages
Requires=online.service
After=online.service

[Service]
Type=oneshot
SyslogIdentifier=pacman
ExecStart=/usr/bin/sudo pacman -Syu --noconfirm
