# systemd/user/patchage.service
[Unit]
Description=jack patch bay

Requisite=wm.target
After=wm.target

Requires=jack.service
After=jack.service

[Service]
Type=simple
SyslogIdentifier=patchage
Slice=apps.slice

ExecStart=/usr/bin/patchage


