# systemd/user/portrait.target
[Unit]
Description=portrait orientation
Conflicts=landscape.target

Requires=yoga-rotate@90.service