# systemd/user/presentation.target
[Unit]
Description=presentation configuration
AllowIsolate=true

Conflicts=console.target desktop.target vncserver.target
