# systemd/user/presync.target
[Unit]
Description=homedir presync hook

StopWhenUnneeded=true

# Shut down all services requiring a synced system.
Conflicts=synced.target