# systemd/user/qsampler.service
[Unit]
Description=linuxsampler frontend

Requisite=wm.target
After=wm.target

Requires=jack.service
After=jack.service

Requires=linuxsampler.service
After=linuxsampler.service

[Service]
Type=simple
SyslogIdentifier=qsampler
Slice=apps.slice

ExecStart=/usr/bin/qsampler


