# systemd/user/reboot.target
[Unit]
Description=reboot system
DefaultDependencies=no
BindsTo=systemd-reboot.service
Before=systemd-reboot.service
