# systemd/user/reboot.target
[Unit]
Description=reboot system

DefaultDependencies=no

Requires=systemd-reboot.service
After=systemd-reboot.service

