# systemd/user/rtags.socket
[Unit]
Description=rtags server socket

[Socket]
ListenStream=%h/.rdm

[Install]
WantedBy=sockets.target
