# systemd/user/shutdown.target
[Unit]
Description=shutdown

DefaultDependencies=no
RefuseManualStart=yes

