# systemd/user/slices.target
[Unit]
Description=user slices

Before=user.target

[Install]
WantedBy=user.target