# systemd/user/sshmux-proxy@.service
[Unit]
Description=ssh connection proxy for %i

BindsTo=sshmux@%i.service
After=sshmux@%i.service

Requires=sshmux@%i.socket

[Service]
Type=simple
Slice=daemons.slice

NonBlocking=true
ExecStart=/usr/lib/systemd/systemd-socket-proxyd %t/ssh-sshmux/%u@%i.home.abdoroig.net.socket
