# systemd/user/stream-stop.service
[Unit]
Description=stop stream player

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user stop stream@*.service

