# systemd/user/sync-fake.service
[Unit]
Description=fake initial sync

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user --no-block start synced.target
