# systemd/user/sync-bye.service
[Unit]
Description=final homedir sync

Conflicts=sync-hi.service sync-fake.service synced.target

Wants=journal-show.service
After=journal-show.service

Requires=homedir-sync@skynet.service
After=homedir-sync@skynet.service

OnFailure=journal-mailer.service


[Service]
Type=oneshot
SyslogIdentifier=sync-bye
ExecStart=/usr/bin/true