# systemd/user/systemd-reboot.service
[Unit]
Description=reboot system

DefaultDependencies=no

[Service]
Type=oneshot

ExecStart=/usr/bin/sleep 1
# NOTE: don't want --force here like the system instance.
ExecStart=/usr/bin/systemctl --no-block reboot

