# systemd/user/test-sleep10.service
[Unit]
Description=sleep for 10 secconds
RemainAfterExit=true

[Service]
Type=simple
ExecStart=/usr/bin/sleep 10
