# systemd/user/test.timer
[Unit]
Description=test timer

[Timer]
OnActiveSec=10
Persistent=true
Unit=test.service
