# systemd/user/test-A.service
[Unit]
Description=test A

[Service]
Type=simple
ExecStart=/usr/bin/echo AAA
