# systemd/user/user.target
[Unit]
Description=user target
AllowIsolate=true
