# systemd/user/vifm-dropdown.service
[Unit]
Description=vifm dropdown

Requisite=graphical.target 
After=graphical.target

[Service]
Type=simple
Slice=apps.slice

ExecStart=/usr/bin/termite --class=vifm-dropdown -e vifm

