# systemd/user/vnc@.target
[Unit]
Description=%i window manager via VNC server

BindsTo=vncserver@2.service
BindsTo=%i@2.service

