# systemd/user/vpn@.target
[Unit]
Description=vpn connection to %i

Requires=openvpn@%i.service
After=openvpn@%i.service
