# systemd/user/wm.slice
[Unit]
Description=window manager slice

Before=slices.target

[Install]
WantedBy=slices.target