# systemd/user/wol@.service
[Unit]
Description=wake on lan %i
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot

ExecStart=/usr/bin/sh -c ' \
    cat /etc/ethers |      \
    grep "^[^# ].*" |      \
    grep "%i" |            \
    cut -d" " -f1 |        \
    xargs wol              \
'
