# systemd/user/work.target
[Unit]
Description=work related servers
StopWhenUnneeded=true

