# systemd/user/xephyr@.target
[Unit]
Description=%i window manager via xephyr

Requisite=xorg@0.service

BindsTo=xephyr@3.service
BindsTo=%i@3.service

