# systemd/user/xephyr@.target
[Unit]
Description=%i window manager via xephyr

Requisite=graphical.target

BindsTo=xephyr@3.service
BindsTo=%i@3.service

