# systemd/user/xephyr-awesome.service
[Unit]
Description=awesome window manager in xephyr

Conflicts=xephyr-openbox.service xephyr-xmonad.service
After=xephyr-openbox.service xephyr-xmonad.service

Requires=xephyr.service
After=xephyr.service

[Service]
Type=simple
SyslogIdentifier=xephyr-awesome
Slice=wm.slice

Environment=DISPLAY=:1
ExecStart=/usr/bin/awesome
KillSignal=SIGINT
TimeoutSec=15

