# systemd/user/xephyr-xmonad.service
[Unit]
Description=xmonad window manager in xephyr

Conflicts=xephyr-awesome.service xephyr-openbox.service

Requires=xephyr.service
After=xephyr.service

[Service]
Type=simple
Slice=wm.slice

Environment=DISPLAY=:1
ExecStart=/usr/bin/xmonad
KillSignal=SIGINT
TimeoutSec=15

