# systemd/user/xorg.target
[Unit]
Description=Plain X.org
Requires=dbus.socket
BindsTo=xorg@0.service
AllowIsolate=true
