# systemd/user/xserver@.target
[Unit]
Description=X server
StopWhenUnneeded=true
RefuseManualStart=true
