# system/systemd/user/dpms-off.service
[Unit]
Description=dpms power down

PartOf=xserver.target
After=xserver.target

[Service]
Type=oneshot
ExecStart=/usr/bin/xset dpms force off
