# system/systemd/user/gnome-session.service
[Unit]
Description=gnome session manager

Conflicts=openbox.service awesome.service

# Need these because this service is ordered before graphical-session
After=xorg.service vncserver.service

BindsTo=graphical-session.target
Before=graphical-session.target

[Service]
# Environment=XDG_SESSION_TYPE=wayland
# Environment=DESKTOP_SESSION=gnome
Environment=XDG_SESSION_TYPE=x11

BusName=org.gnome.SessionManager
ExecStart=/usr/bin/gnome-session -a /dev/null
