# system/systemd/user/i3lock.service
[Unit]
Description=i3lock screen locker

Requisite=graphical-session.target
After=graphical-session.target

Requires=lock-pre.target
After=lock-pre.target

# so that when lock.target stops, i3lock goes down
PartOf=lock.target

[Service]
Type=simple
ExecStart=/usr/bin/i3lock --nofork --color=000000


[Install]
WantedBy=lock.target
