# system/systemd/user/lock-session.target
[Unit]
Description=lock session

After=disconnect.target

DefaultDependencies=no
RefuseManualStart=yes
StopWhenUnneeded=yes
