# system/systemd/user/mako.service
[Unit]
Description=Mako notification daemon
Documentation=man:mako(1)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
Slice=wm.slice

BusName=org.freedesktop.Notifications
ExecStart=/usr/bin/mako

[Install]
WantedBy=default.target

