# system/systemd/user/openbox.service
[Unit]
Description=openbox window manager

Conflicts=awesome.service

# Need these because this service is ordered before graphical-session
After=xorg.service vncserver.service xrdb.service

BindsTo=graphical-session.target
Before=graphical-session.target

[Service]
Type=simple
Slice=wm.slice

ExecStart=/usr/bin/openbox
KillSignal=SIGINT
TimeoutSec=15
