# system/systemd/user/sway-session.target
[Unit]
Description=sway wm session
AllowIsolate=true

BindsTo=sway.service
After=sway.service
