# system/systemd/user/xautolock-inhibit.service
[Unit]
Description=temporarily inhibit autolock

StopWhenUnneeded=true

Requisite=graphical-session.target
After=graphical-session.target

Requisite=xautolock.service
After=xautolock.service

[Service]
Type=oneshot
ExecStart=/usr/bin/xautolock -disable
ExecStop=/usr/bin/xautolock -enable

