# user/machine-ephimeral@.service.d/ephimeral.conf
[Unit]
Description=ephimeral qemu vm for %i

Requires=machine-snapshot@%i.service
After=machine-snapshot@%i.service

[Service]
Environment=ROOT=/home/abdo/virt/tmp/%i.qcow2