# user/gvfs-daemon.service.d/disable-fuse.conf
[Service]
Environment=GVFS_DISABLE_FUSE=1
