# user/a2jmidid.service
[Unit]
Description=alsa to jack midi bridge

Requires=jack.service
After=jack.service

Requires=audio.slice
After=audio.slice

[Service]
Type=simple
SyslogIdentifier=a2jmidid
Slice=audio.slice

# realtime scheduling
CPUSchedulingPolicy=rr
CPUSchedulingPriority=60

ExecStart=/usr/bin/a2jmidid -u -e
