# user/ac.target
[Unit]
Description=ac power
Conflicts=battery.target
