# user/adb.service
[Unit]
Description=android debug bridge daemon

[Service]
Type=forking
ExecStart=/opt/android-sdk/platform-tools/adb -P 5038 start-server
