# user/apps.slice
[Unit]
Description=Applications slice

Requisite=xorg.target
StopWhenUnneeded=true
