# user/apps.slice
[Unit]
Description=user applications

[Install]
WantedBy=wm.target
