# user/apps.slice
[Unit]
Description=applications slice

Requisite=desktop.target
After=desktop.target

StopWhenUnneeded=true
