# user/apps.slice
[Unit]
Description=applications

Requisite=desktop.target
After=desktop.target
