# user/apps.slice
[Unit]
Description=user applications

Before=slices.target

[Install]
WanedBy=slices.target