# user/apps.slice
[Unit]
Description=user applications

Before=slices.target

[Install]
WantedBy=wm.target