# user/apps.service
[Unit]
Description=manage apps.slice

Requires=desktop.target apps.slice
After=desktop.target apps.slice

[Service]
Type=oneshot
RemainAfterExit=true

ExecStart=/usr/bin/systemctl --user start apps.slice
ExecStop=/usr/bin/systemctl --user kill apps.slice
