# user/apps.slice
[Unit]
Description=user applications

[Slice]
MemoryHigh=70%
MemoryMax=80%

