# user/async-default-base.service
[Unit]
Description=sync mail to default host

Requires=async-default-mounted.service
After=async-default-mounted.service

[Service]
Type=oneshot
ExecStart=/usr/bin/async -b sync --terse default -d etc,priv,var
