# user/async-needed@.service
[Unit]
Description=synchronize homedir to %i

Requires=mount-priv.service
After=mount-priv.service

Requires=network-online.service
After=network-online.service

[Service]
Type=oneshot
SyslogIdentifier=async-%i
KillSignal=SIGINT
ExecStart=/usr/bin/async -b sync %i --terse --needed 
