# user/async@.service
[Unit]
Description=synchronize with %i

Requires=network-online.service
After=network-online.service

[Service]
Type=oneshot
SyslogIdentifier=async
KillSignal=SIGINT
ExecStart=/usr/bin/async -b sync %i --terse
