# user/async-r2d2.service
[Unit]
Description=sync directories to r2d2

# BindsTo=sys-subsystem-net-devices-r2d2.device
# After=sys-subsystem-net-devices-r2d2.device

[Service]
Type=oneshot
SyslogIdentifier=async-r2d2
KillSignal=SIGINT
ExecStart=/usr/bin/async -b sync --terse r2d2