# user/async-default-work.service
[Unit]
Description=sync mail to default host

Requires=async-default-mounted.service
After=async-default-mounted.service

# so we do not get two unison syncs at the same time
After=async-default-base.service

[Service]
Type=oneshot
ExecStart=/usr/bin/async -b sync --terse default -d worki