# user/audio.target
[Unit]
Description=Audio servers
AllowIsolate=true

