# user/audio.slice
[Unit]
Description=audio servers

Before=slices.target

[Install]
WantedBy=slices.target
