# user/audio.target
[Unit]
Description=Audio servers
Requires=dbus.socket
AllowIsolate=true

[Install]
WantedBy=default.target
