# user/audio.slice
[Unit]
Description=audio servers

Before=slices.target

[Install]
WanedBy=slices.target
