# user/awesome.target
[Unit]
Description=Awesome Window Manager
Requires=dbus.socket
BindsTo=awesome.service
AllowIsolate=true
