# user/awesome.target
[Unit]
Description=awesome desktop environment
Conflicts=console.target

BindsTo=awesome.service
After=awesome.service
