# user/awesome@.service
[Unit]
Description=Awesome Window Manager
After=xorg.target
StopWhenUnneeded=true

Requires=xorg.target
Requires=dbus.socket

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/conf/xenvironment
Environment=DISPLAY=:%i
ExecStart=/home/abdo/.config/systemd/scripts/awesome start
KillSignal=SIGINT
TimeoutSec=15

[Install]
WantedBy=awesome.target
