# user/awesome@.service
[Unit]
Description=awesome window manager at display %i

After=xorg@%i.service vncserver@%i.service xephyr@%i.service
After=compton.service xrandr-auto.service

Requisite=xserver@%i.target
After=xserver@%i.target

Before=wm.target

[Service]
Type=simple
SyslogIdentifier=awesome
Slice=wm.slice

Environment=DISPLAY=:%i
ExecStart=/usr/bin/awesome
KillSignal=SIGINT
TimeoutSec=15
