# user/awesome.service

[Unit]
Description=The Awesome Window Manager
After=xorg.target
Requires=xorg.target
Requires=dbus.socket

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/conf/xenvironment
# ExecStart=/usr/bin/awesome
ExecStart=/home/abdo/.config/systemd/scripts/awesome start
ExecStop=/home/abdo/.config/systemd/scripts/awesome stop

[Install]
WantedBy=awesome.target
