# user/awesome.service
[Unit]
Description=Awesome Window Manager
After=xorg.target
StopWhenUnneeded=true
Requires=xorg.target
Requires=dbus.socket
Conflicts=xmonad.service openbox.service

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/conf/xenvironment
# ExecStart=/usr/bin/awesome
ExecStart=/home/abdo/.config/systemd/scripts/awesome start
# ExecStop=/home/abdo/.config/systemd/scripts/awesome stop

[Install]
WantedBy=awesome.target
