# user/awesome.target

[Unit]
Description=Awesome Window Manager
Requires=dbus.socket
AllowIsolate=true
