# user/awesome@.service
[Unit]
Description=awesome window manager at display %i

After=xserver@%i.target
After=xorg@%i.service vncserver@%i.service x11vnc@%i.service xephyr@%i.service
After=compton.service xrandr-auto.service
Before=desktop.target

Requisite=xserver@%i.target

Requires=wm.target
PartOf=wm.target

[Service]
Type=simple
SyslogIdentifier=awesome
Slice=wm.slice

Environment=DISPLAY=:%i
ExecStart=/usr/bin/awesome
KillSignal=SIGINT
TimeoutSec=15
