# user/awesome-session.target
[Unit]
Description=awesome wm session
AllowIsolate=true

BindsTo=graphical-session.target
Before=graphical-session.target

BindsTo=xorg.service
After=xorg.service

BindsTo=awesome.service
After=awesome.service

Wants=compton.service redshift.service

Conflicts=tint2.service
