# user/awesome.service

[Unit]
Description=The Awesome Window Manager
After=xorg.target
Requires=dbus.socket

[Service]
EnvironmentFile=%h/.config/systemd/conf/xenvironment
ExecStart=/usr/bin/awesome
Restart=on-failure

[Install]
WantedBy=awesome.target
