# user/bamboo.service
[Unit]
Description=Setup bamboo tablet coordinates
After=xorg.target
ConditionHost=grothendieck

Requisite=xorg.target

[Service]
Type=oneshot
ExecStart=/home/abdo/.config/systemd/scripts/bamboo.sh

[Install]
WantedBy=desktop.target
