# user/battery.target
[Unit]
Description=battery power
Conflicts=ac.target
