# user/btpd.service
[Unit]
Description=Bit Torrent Protocol Daemon

[Service]
Type=simple
ExecStart=/usr/bin/btpd --no-daemon --bw-in 300 --bw-out 30 -p 6881
