# user/bye.target
[Unit]
Description=trigger final actions
StopWhenUnneeded=true

# If a target fails to start notify by mail.
Requires=journal-mailer.service
After=journal-mailer.service

# refresh cached passphrase now
Requires=gpg-agent-auth.service

After=final-synced-data.target final-synced-mail.target
Before=final-shutdown.service
Requires=final-synced-data.target final-synced-mail.target final-shutdown.service
