# user/bye.target
[Unit]
Description=trigger final actions
StopWhenUnneeded=true

After=final-synced-data.target final-synced-mail.target
Before=final-shutdown.service
Requires=final-synced-data.target final-synced-mail.target final-shutdown.service
