# user/bye.target
[Unit]
Description=trigger final actions and poweroff
StopWhenUnneeded=true

OnFailure=journal-mailer.service

# refresh cached passphrase now
Requires=gpg-agent-auth.service

Requires=final-synced-data.target
After=final-synced-data.target

Requires=delayed-poweroff.service
Before=delayed-poweroff.service


