# user/bye.target
[Unit]
Description=trigger final actions and poweroff
StopWhenUnneeded=true

OnFailure=journal-mailer.service

Requires=final-synced-all.target
After=final-synced-all.target

Requires=final-synced-poweroff.service
Before=final-synced-poweroff.service


