# user/bye.target
[Unit]
Description=trigger final actions and poweroff
StopWhenUnneeded=true

OnFailure=journal-mailer.service

# refresh cached passphrase now
Requires=gpg-agent-auth.service

Requires=final-synced-all.target
After=final-synced-all.target

Requires=delayed-poweroff.target
Before=delayed-poweroff.target


