# user/bye.target
[Unit]
Description=trigger final actions

StopWhenUnneeded=true
AllowIsolate=true
