# user/byez.target
[Unit]
Description=trigger final actions and suspend
StopWhenUnneeded=true

OnFailure=journal-mailer.service

# refresh cached passphrase now
Requires=gpg-agent-auth.service

Requires=final-synced-all.target
After=final-synced-all.target

Requires=delayed-suspend.target
Before=delayed-suspend.target

