# user/byez.target
[Unit]
Description=trigger final actions and sleep
StopWhenUnneeded=true

# If it fails, email me a piece of journal
OnFailure=journal-mailer.service

# refresh cached passphrase now
Requires=gpg-agent-auth.service

After=final-synced-data.target final-synced-mail.target
Before=final-suspend.service
Requires=final-synced-data.target final-synced-mail.target final-suspend.service
