# user/calibre-server.service
[Unit]
Description=Calibre Server

[Service]
Type=simple
ExecStart=/usr/bin/calibre-server --with-library=/home/abdo/lib
