# user/calibre.service
[Unit]
Description=calibre
After=xorg.target

Requisite=xorg.target
Requires=calibre-unlock.service

[Service]
Type=simple
SyslogIdentifier=calibre

ExecStart=/home/abdo/.config/systemd/scripts/calibre.sh start
ExecStop=/home/abdo/.config/systemd/scripts/calibre.sh stop

