# user/calibre.service
[Unit]
Description=calibre
After=xorg.target

Requisite=desktop.target 
After=desktop.target

Requires=calibre-db.service
After=calibre-db.service

[Service]
Type=simple
SyslogIdentifier=calibre
Slice=apps.slice

ExecStart=/usr/bin/calibre --with-library "$AB2_CALIBRE_DIR"


