# user/calibre-db.service
[Unit]
Description=unlock calibre database

StopWhenUnneeded=true
ConditionPathExists=${AB2_CALIBRE_DIR}

Requisite=initial-synced-all.service
After=initial-synced-all.service

[Service]
Type=oneshot
RemainAfterExit=true
SyslogIdentifier=calibre-db

ExecStart=/usr/bin/true
ExecStop=/usr/bin/systemctl --user start commit-calibre.service
