# user/calibre-commit.service
[Unit]
Description=commit calibre library

Requisite=initial-synced-data.target
After=initial-synced-data.target

Requires=calibre-stop.service
After=calibre-stop.service

[Service]
Type=oneshot

SyslogIdentifier=calibre-commit
ExecStart=/home/abdo/.config/systemd/scripts/git-commit.sh "$AB2_CALIBRE_DIR"
