# user/calibre-commit.service
[Unit]
Description=commit calibre library

Conflicts=calibre.service calibre-db.service
Before=calibre.service calibre-db.service

[Service]
Type=oneshot

SyslogIdentifier=calibre-commit
ExecStart=/home/abdo/.config/systemd/scripts/calibre.sh commit
