# user/calibre-server.service
[Unit]
Description=Calibre Server

ConditionPathExists=/home/%u/lib

Requires=git-use@lib.service
After=git-use@lib.service

[Service]
Type=simple
ExecStart=/usr/bin/calibre-server --with-library "${AB2_CALIBRE_DIR}"
