# user/calibre.service
[Unit]
Description=calibre
After=xorg.target

Requisite=desktop.target 
Requires=calibre-db.service
After=calibre-db.service

[Service]
Type=simple
SyslogIdentifier=calibre

ExecStart=/home/abdo/.config/systemd/scripts/calibre.sh start
ExecStop=/home/abdo/.config/systemd/scripts/calibre.sh stop

