# user/calibre-db.service
[Unit]
Description=unlock calibre database
StopWhenUnneeded=true

[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=calibre-db
ExecStart=/home/abdo/.config/systemd/scripts/calibre.sh unlock
ExecStop=/home/abdo/.config/systemd/scripts/calibre.sh commit
