# user/chat-away.service
[Unit]
Description=chat away state
ConditionPathExists=%t/emacs/chat

[Service]
Type=oneshot
RemainAfterExit=yes

StandardOutput=null
ExecStart=/home/abdo/bin/emacs-run -q -s chat.service -k chat '(abdo-chat-away)'
ExecStop=/home/abdo/bin/emacs-run -q -s chat.service -k chat '(abdo-chat-back)'
