# user/chat-db.service
[Unit]
Description=chat log database

StopWhenUnneeded=true

[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=chat-db
ExecStart=/usr/bin/true
ExecStop=/home/abdo/.config/systemd/scripts/git-commit.sh "$AB2_VAR_DIR/chat/emacs"

