# user/chat-away.service
[Unit]
Description=chat away state
ConditionPathExists=/tmp/emacs%U/chat
StopWhenUnneeded=true

[Service]
Type=oneshot
RemainAfterExit=yes

StandardOutput=null
ExecStart=/usr/bin/emacsclient -s chat -e '(abdo-chat-away)'
ExecStop=/usr/bin/emacsclient -s chat -e '(abdo-chat-back)'
