# user/chat.service
[Unit]
Description=chat

AssertPathExists=/home/%u/priv

Requisite=graphical.target
After=graphical.target

Requires=git-use@var-chat-emacs.service
After=git-use@var-chat-emacs.service

[Service]
Type=simple
Slice=apps.slice

ExecStart=/usr/bin/emacsclient -s "%t/emacs/server" -c --frame-parameters='((title . "chat"))' --eval '(ab2/frame-layout "@chat")'
KillSignal=SIGINT
