# user/chat-commit.service
[Unit]
Description=commit chat logs

Requisite=initial-synced-data.target
After=initial-synced-data.target

Requires=chat-stop.service
After=chat-stop.service

[Service]
Type=oneshot

SyslogIdentifier=chat-commit
ExecStart=/home/abdo/.config/systemd/scripts/git-commit.sh "$AB2_VAR_DIR/chat/emacs"
