# user/chat.service
[Unit]
Description=chat

Requisite=desktop.target
After=desktop.target

Requires=chat-db.service
After=calibre-db.service

[Service]
Type=simple
SyslogIdentifier=chat

ExecStart=/home/abdo/.config/systemd/scripts/chat.sh start
ExecStop=/home/abdo/.config/systemd/scripts/chat.sh stop

