# user/chat-commit.service
[Unit]
Description=commit chat logs

Requires=chat-stop.service

[Service]
Type=oneshot

SyslogIdentifier=chat-commit
ExecStart=/home/abdo/.config/systemd/scripts/chat.sh commit
