# user/check-mail.timer
[Timer]
OnActiveSec=5
OnUnitInactiveSec=2min

Unit=check-mail.service
