# user/check-mail.service
[Unit]
Description=check for new mail

[Service]
Type=oneshot
ExecStart=/home/abdo/.config/systemd/scripts/check-mail.py
