# user/claudia.service
[Unit]
Description=jack patch bay and ladish client

PartOf=graphical-session.target
After=graphical-session.target

Requires=jack.service
After=jack.service

[Service]
Type=simple
Slice=apps.slice

ExecStart=/usr/bin/claudia


