# user/clean-build.service
[Unit]
Description=clean build files

PartOf=housekeeping.target

AssertPathExists=%h/build/Makefile

[Service]
Type=oneshot
Slice=workers.slice

# NOTE: Do not use environment variables here. If they are empty, it can mess things up royally.
ExecStart=/usr/bin/make --silent -C "%h/build" random-expire
