# user/commit-chat.service
[Unit]
Description=commit chat logs

ConditionPathExists=${AB2_VAR_DIR}/chat/emacs

Conflicts=chat.service chat-db.service

[Service]
Type=oneshot

SyslogIdentifier=commit-chat
ExecStart=/home/abdo/.config/systemd/scripts/git-commit.sh "${AB2_VAR_DIR}/chat/emacs"

