# user/commit-orgmode.service
[Unit]
Description=commit orgmode

ConditionPathExists=/home/%u/work/wiki

Conflicts=orgmode-db.service orgmode.service 

[Service]
Type=oneshot

SyslogIdentifier=commit-orgmode
ExecStart=/home/abdo/.config/systemd/scripts/git-commit.sh "${AB2_WIKI_DIR}"
