# user/commit-calibre.service
[Unit]
Description=commit calibre library

ConditionPathExists=/home/%u/lib

Conflicts=calibre.service calibre-server.service calibre-db.service

[Service]
Type=oneshot

SyslogIdentifier=commit-calibre
ExecStart=/home/abdo/.config/systemd/scripts/git-commit.sh "${AB2_CALIBRE_DIR}"

# NOTE: if we change calibre config in the database, I need to trigger a commit manually
#       to run the git hooks and dump the database. Otherwise the working tree appears clean
#       and nothing happens.