# user/compositor.target
[Unit]
Description=compositor

Requires=graphical.target

[Install]
WantedBy=desktop.target
WantedBy=laptop.target
WantedBy=tablet.target

