# user/compton.service
[Unit]
Description=Compton X compositor
After=xorg.service
StopWhenUnneeded=true
Requires=xorg.service

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/conf/xenvironment
ExecStart=/usr/bin/compton

[Install]
WantedBy=xorg.target
