# user/compton.service
[Unit]
Description=Compton X compositor

Requisite=xorg.target
After=xorg.target

[Service]
Type=simple
SyslogIdentifier=compton
Slice=wm.slice

ExecStart=/usr/bin/compton
