# user/compton.service
[Unit]
Description=Compton X compositor
After=xorg.target

Requisite=xorg.target

[Service]
Type=simple
SyslogIdentifier=compton
Slice=wm.slice

ExecStart=/usr/bin/compton
