# user/compositor.target
[Unit]
Description=compositor

Requires=user.target
After=user.target

Requires=graphical.target

[Install]
WantedBy=desktop.target
WantedBy=laptop.target
WantedBy=tablet.target

