# user/connman-ui.service
[Unit]
Description=Connman UI
ConditionHost=galois

Requisite=xorg.target
Requires=gvfsd.service gconf.service
After=gvfsd.service gconf.service

[Service]
Type=simple
ExecStart=/usr/bin/connman-ui-gtk
