# user/connman-ui.service
[Unit]
Description=Connman UI
ConditionHost=galois

[Service]
Type=simple
ExecStart=/usr/bin/connman-ui-gtk
