# user/console.target
[Unit]
Description=console configuration
AllowIsolate=true

