# user/console.target
[Unit]
Description=tty console
StopWhenUnneeded=true

Conflicts=xorg.target

AllowIsolate=true
