# user/console.target
[Unit]
Description=tty console
Conflicts=xorg.target

AllowIsolate=true
