# user/console.target
[Unit]
Description=tty console

StopWhenUnneeded=true
AllowIsolate=true

Conflicts=xorg.target


