# user/console.target
[Unit]
Description=console configuration
AllowIsolate=true

Conflicts=graphical.target 
Conflicts=laptop.target desktop.target tablet.target


