# user/contacts.service
[Unit]
Description=fetch google contacts

[Service]
Type=oneshot
ExecStart=/usr/bin/gcont -d -p %h/mail/contacts
