# user/daemons.target
[Unit]
Description=user daemons

[Install]
WantedBy=default.target
