# user/daemons.slice
[Unit]
Description=user daemons

Before=slices.target

[Install]
WanedBy=slices.target