# user/daemons.slice
[Unit]
Description=user daemons

Before=slices.target

[Install]
WantedBy=slices.target