# user/daily.timer
[Unit]
Description=Daily Timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d
Unit=daily.target

[Install]
WantedBy=basic.target
