# user/default.target
[Unit]
Description=Default target
AllowIsolate=true
