# user/delay@.service
[Unit]
Description=delay by %i seconds

[Service]
Type=oneshot
SyslogIdentifier=delay
ExecStart=/usr/bin/sleep %i

