# user/delayed-suspend.service
[Unit]
Description=suspend system in 20 seconds
DefaultDependencies=no

BindsTo=systemd-suspend.service
Before=systemd-suspend.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sleep 20