# user/delayed-reboot.service
[Unit]
Description=reboot system in 20 seconds
DefaultDependencies=no

BindsTo=systemd-reboot.service
Before=systemd-reboot.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sleep 20