# user/delicious.service
[Unit]
Description=fetch delicious bookmarks
Requires=network-online.service
After=network-online.service

ConditionPathExists=/home/%u/.netrc

[Service]
Type=oneshot
SyslogIdentifier=delicious
ExecStart=/home/abdo/.config/systemd/scripts/fetch-delicious.sh 
