# user/desktop.target
[Unit]
Description=desktop environment

StopWhenUnneeded=true
AllowIsolate=true

Conflicts=console.target

Requires=xorg.target
After=xorg.target

