# user/desktop.target
[Unit]
Description=desktop configuration
AllowIsolate=true

Conflicts=console.target laptop.target tablet.target vncserver.target

