# user/desktop.target
[Unit]
Description=desktop environment
AllowIsolate=true

Conflicts=console.target

Requires=graphical.target
After=graphical.target

