# user/desktop.target
[Unit]
Description=Desktop Environment
Requires=user.target
After=xorg.target user.target
StopWhenUnneeded=true

Requisite=xorg.target

AllowIsolate=true
