# user/desktop.target
[Unit]
Description=Desktop Environment
Requires=dbus.socket
Requires=xorg.target
StopWhenUnneeded=true
AllowIsolate=true
