# user/desktop.target
[Unit]
Description=desktop environment

Requires=xorg.target
After=xorg.target

StopWhenUnneeded=true
AllowIsolate=true
