# user/desktop.target
[Unit]
Description=Desktop Environment
After=xorg.target
StopWhenUnneeded=true

Requires=dbus.socket
Requisite=xorg.target

AllowIsolate=true
