# user/dictd.service
[Unit]
Description=dictd dictionary server daemon

ConditionPathExists=%h/etc/dictd/dictd.conf
ConditionHost=!quark
ConditionHost=!skynet

[Service]
Type=forking
Slice=daemons.slice

SuccessExitStatus=143
ExecStart=/usr/sbin/dictd -c %h/etc/dictd/dictd.conf --pid-file %t/dictd.pid

[Install]
WantedBy=daemons.target