# user/dictd.service
[Unit]
Description=Dictd Dictionary Server Daemon

[Service]
Type=forking

SuccessExitStatus=143
ExecStart=/usr/sbin/dictd -c %h/etc/dictd/dictd.conf --pid-file %t/dictd.pid
