# user/disconnect.target
[Unit]
Description=close active connections

DefaultDependencies=no
RefuseManualStart=yes

