# user/disconnect.target
[Unit]
Description=close active network services

Before=sleep.target shutdown.target

[Install]
WantedBy=sleep.target
WantedBy=shutdown.target
