# user/disconnect@.target
[Unit]
Description=close network links to %i

StopWhenUnneeded=yes

