# user/dpms-off.service
[Unit]
Description=dpms power down

Requisite=graphical.target
After=graphical.target

[Service]
Type=oneshot
SyslogIdentifier=dpms-off
ExecStart=/usr/bin/xset dpms force off