# user/dropdown.slice
[Unit]
Description=dropdown applications

[Install]
WantedBy=wm.target
