# user/dropbox.service
[Unit]
Description=dropbox daemon

Requires=dbus.socket

[Service]
Type=simple
ExecStart=/usr/bin/dropboxd
