# user/dropdown.slice
[Unit]
Description=dropdown applications

Requisite=desktop.target
After=desktop.target
