# user/dropbox.service
[Unit]
Description=dropbox daemon

Requisite=xorg@0.service
Requires=dbus.socket

[Service]
Type=simple
ExecStart=/usr/bin/dropboxd
