# user/dropdown.slice
[Unit]
Description=dropdown applications

Before=slices.target

[Install]
WantedBy=slices.target