# user/dropbox.service
[Unit]
Description=dropbox daemon

After=xorg.target
Requisite=xorg.target
Requires=dbus.socket

[Service]
Type=simple
ExecStart=/usr/bin/dropboxd
