# user/dropdown.slice
[Unit]
Description=dropdown applications slice

Requisite=desktop.target
After=desktop.target
