# user/dwb-db.service
[Unit]
Description=dwb cookie database
StopWhenUnneeded=true

Requires=mount-priv.service
After=mount-priv.service

# dwb cookies are under ~/priv
Requisite=initial-synced-all.service
After=initial-synced-all.service


[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=dwb-db
ExecStart=/usr/bin/true
ExecStop=/usr/bin/systemctl --user start commit-dwb.service
