# user/dwb-commit.service
[Unit]
Description=commit dwb cookies

Requisite=initial-synced-data.target
After=initial-synced-data.target

Requires=dwb-stop.service
After=dwb-stop.service

[Service]
Type=oneshot

SyslogIdentifier=dwb-commit
ExecStart=/home/abdo/.config/systemd/scripts/git-commit.sh "$AB2_PRIV_DIR/var/dwb/default"
