# user/dwb-stop.service
[Unit]
Description=dwb stop

Conflicts=dwb-db.service

[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=dwb-stop
ExecStart=/home/abdo/.config/systemd/scripts/dwb.sh stop
