# user/dwb@.service
[Unit]
Description=dwb %i

Requisite=wm.target 
After=wm.target

Requires=git-use@priv-var-dwb-default.service
After=git-use@priv-var-dwb-default.service

[Service]
Type=simple
SyslogIdentifier=dwb
Slice=apps.slice

# Start new instance of dwb without dbus
ExecStart=/usr/bin/dwb -n
StandardOutput=null

