# user/ec2-stop@.service
[Unit]
Description=start ec2 instance of %i

Requires=network-online.target
After=network-online.target

Requires=mount-priv.service
After=mount-priv.target

[Service]
Type=oneshot
SyslogIdentifier=ec2-stop

ExecStart=/home/abdo/.config/systemd/scripts/ec2-host.sh stop %i
