# user/emacs-save-quit.service
[Unit]
Description=emacs save and quit

[Service]
Type=simple
SyslogIdentifier=emacs-stop
ExecStart=/home/abdo/.config/systemd/scripts/emacs.sh stop

