# user/emacs-save.service
[Unit]
Description=save emacs buffers

# So that emacs.service is stoped after this. The stopping order gets reversed.
After=emacs.service

[Service]
Type=oneshot

# NOTE: If emacs.service is not active, always succeed since there is nothing to be done.
#       It would be nice to write the dependency on emacs.service directly in systemd!
ExecStart=/home/abdo/bin/emacs-run -q -s emacs.service -k server '(save-some-buffers t)'

