# user/emacs.socket
[Unit]
Description=emacs daemon socket

[Socket]
ListenStream=%t/emacs/server
Service=emacs-proxy.service

[Install]
WantedBy=sockets.target
