# user/emacsclient@.service
[Unit]
Description=emacs client %i

Requisite=wm.target 
After=wm.target

[Service]
Type=simple
SyslogIdentifier=emacsclient
Slice=apps.slice

ExecStart=/usr/bin/emacsclient -s "${XDG_RUNTIME_DIR}/emacs/server" -c


