# user/emacsclient@.service
[Unit]
Description=emacs client %i

Requisite=wm.target 
After=wm.target

[Service]
Type=simple
SyslogIdentifier=emacsclient
Slice=apps.slice

ExecStart=/usr/bin/emacsclient -s /run/user/1000/emacs/server -c


