# user/emacs.service
[Unit]
Description=emacs daemon

# need a dependency on xorg. otherwise, when X quits emacs crashes.
# emacsclient says something about a long standing GTK bug.
Requisite=graphical.target
After=graphical.target

[Service]
Type=forking
Slice=daemons.slice

SyslogIdentifier=emacs

ExecStart=/usr/bin/emacs --daemon
SuccessExitStatus=0 15

# Kill with a SIGINT. This calls kill-emacs. No buffer saving is done.
KillSignal=SIGINT

[Install]
WantedBy=wm.target
