# user/ethereum.service
[Unit]
Description=ethereum daemon
After=network.target

Requires=mount-priv.service
After=mount-priv.service

Requires=git-use@priv-wallet.service
After=git-use@priv-wallet.service

Requires=synced.target
After=synced.target

[Service]
Type=simple
Slice=daemons.slice

ExecStart=/usr/bin/geth --support-dao-fork

# explicitly remove IPC socket
ExecStopPost=/usr/bin/rm -f %h/.ethereum/geth.ipc
