# user/fetch-news.target
[Unit]
Description=fetch news

Requires=initial-synced-data.target
After=initial-synced-data.target
After=journal-mailer.service

Requires=mail-db.service
After=mail-db.service

Requires=mutag.service rss2maildir.service
After=mutag.service rss2maildir.service
