# user/fetch-skynet.target
[Unit]

ConditionHost=skynet

Description=fetch mail and news

# stopping of git-use@mail propagates to this target
PartOf=git-use@mail.service
Before=git-use@mail.service

Requires=offlineimap.service rss2maildir.service
Before=offlineimap.service rss2maildir.service
