# user/fetch-news.target
[Unit]
Description=fetch news

# stopping of git-use@mail propagates to this target
PartOf=git-use@mail.service
Before=git-use@mail.service

Requires=rss2maildir.service
Before=rss2maildir.service

Requires=mutag-autotag.service mutag-index.service
Before=mutag-autotag.service mutag-index.service