# user/fetch-mail.target
[Unit]
Description=fetch mail
Requires=dbus.socket
StopWhenUnneeded=true
AllowIsolate=true
Wants=offlineimap.service mutag.service

