# user/fetch-mail.timer
[Timer]
OnActiveSec=5
OnUnitInactiveSec=3min

Unit=fetch-mail.service
