# user/fetch-contacts.target
[Unit]
Description=fetch contacts

# stopping of git-use@mail propagates to this target
PartOf=git-use@mail.service
Before=git-use@mail.service

Requires=gcont-sync.service
Before=gcont-sync.service
