# user/fetch-mail.target
[Unit]
Description=fetch mail
StopWhenUnneeded=true

Requires=online.service
After=online.service

Wants=offlineimap.service mutag.service

Requires=mail-db.service
After=mail-db.service
