# user/fetch-skynet.target
[Unit]

ConditionHost=skynet

Description=fetch mail and news
StopWhenUnneeded=true

Requires=git-use@mail.service
After=git-use@mail.service

Requires=offlineimap.service
After=offlineimap.service

Requires=rss2maildir.service
After=rss2maildir.service

