# user/fetch-mail.target
[Unit]
Description=fetch mail
StopWhenUnneeded=true

Requires=offlineimap.service
After=offlineimap.service

Requires=mutag-autotag.service mutag-index.service
After=mutag-autotag.service mutag-index.service