# user/fetch-contacts.target
[Unit]
Description=fetch contacts

Requires=online.service
After=online.service

StopWhenUnneeded=true
AllowIsolate=true

Wants=gcont.service

Requires=mail-unlock.service
After=mail-unlock.service
