# user/fetch-contacts.target
[Unit]
Description=fetch contacts
StopWhenUnneeded=true

Requires=online.service
After=online.service

Wants=gcont-sync.service

Requires=mail-db.service
After=mail-db.service
