# user/fetch.target
[Unit]
Description=fetch mail and news

Requires=initial-synced-data.target
After=initial-synced-data.target
After=journal-mailer.service

Requires=mail-db.service
After=mail-db.service

Requires=gcal-sync.service gcont-sync.service mutag.service offlineimap.service rss2maildir.service
After=gcal-sync.service gcont-sync.service mutag.service offlineimap.service rss2maildir.service
