# user/fetch.target
[Unit]
Description=fetch mail and news
StopWhenUnneeded=true

Requires=mail-db.service
After=mail-db.service

Requires=gcal-sync.service gcont-sync.service mutag-autotag.service offlineimap.service rss2maildir.service
After=gcal-sync.service gcont-sync.service mutag-autotag.service offlineimap.service rss2maildir.service
