# user/fetch-news.target
[Unit]
Description=fetch news

Requires=online.service
After=online.service

StopWhenUnneeded=true
AllowIsolate=true

Wants=rss2maildir.service mutag.service

Requires=mail-db.service
After=mail-db.service
