# user/fetch-news.timer
[Unit]
Description=periodic news fetching

[Timer]
OnActiveSec=5
OnUnitInactiveSec=2h
Unit=fetch-news.target

[Install]
WantedBy=timers.target