# user/fetch-skynet.timer
[Unit]
ConditionHost=skynet

[Timer]

OnActiveSec=5
OnUnitInactiveSec=6h
Unit=fetch-mail.target

[Install]
WantedBy=timers.target