# user/fetch.target
[Unit]
Description=fetch mail and news

StopWhenUnneeded=true
AllowIsolate=true
Wants=offlineimap.service rss2maildir.service gcont.service mutag.service

