# user/fetch-contacts.target
[Unit]
Description=fetch contacts
StopWhenUnneeded=true

Requires=git-use@mail.service
After=git-use@mail.service

Requires=gcont-sync.service
After=gcont-sync.service
