# user/final-synced-data.target
[Unit]
Description=final data sync

# Can't conflict here, because final-synced-mail 
# requires initial-synced.data.
# Conflicts=initial-synced-data.target

Requires=async@default.service
After=async@default.service

