# user/final-sleep.service
[Unit]
Description=final sleep

Requires=final-synced-data.target
After=final-synced-data.target

[Service]
Type=oneshot
SyslogIdentifier=sleep
ExecStart=/usr/bin/sh -c "sleep 10; sudo systemctl suspend"
