# user/final-shutdown.service
[Unit]
Description=final shutdown

# prevent shutdown if final-sync.target fails
Requisite=final-synced-data.target
After=final-synced-data.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl poweroff
