# user/final-synced-data.target
[Unit]
Description=final data sync

Requires=final-synced-mail.target
After=final-synced-mail.target

Requires=async@default.service
After=async@default.service

