# user/final-shutdown.service
[Unit]
Description=final shutdown

Requires=final-synced-data.target
After=final-synced-data.target

[Service]
Type=oneshot
SyslogIdentifier=poweroff
ExecStart=/usr/bin/sh -c "sleep 10; sudo systemctl poweroff"
