# user/final-synced-all.target
[Unit]
Description=final homedir sync

Conflicts=initial-synced-all.target

Requires=async-default-all.service
After=async-default-all.service

