# user/final-synced-mail.target
[Unit]
Description=final mail sync

# If a target fails to start notify by mail.
Requires=journal-mailer.service
After=journal-mailer.service

Requires=initial-synced-data.target
After=initial-synced-data.target

Requires=fetch.target
After=fetch.target
