# user/firefox@.service
[Unit]
Description=firefox %i

Requisite=wm.target 
After=wm.target

[Service]
Type=simple
SyslogIdentifier=firefox
Slice=apps.slice

ExecStart=/usr/bin/firefox
StandardOutput=null


