# user/fluidsynth.service
[Unit]
Description=fluidsynth daemon

Requires=jack.service
After=jack.service

Requires=audio.slice
After=audio.slice

[Service]
Type=simple
ExecStart=/usr/bin/fluidsynth --server --no-shell --audio-driver jack --midi-driver jack --sample-rate 48000 --gain 3.0

# realtime scheduling
CPUSchedulingPolicy=rr
CPUSchedulingPriority=60

Slice=audio.slice

[Install]
WantedBy=multi-user.target
