# user/gcal-sync.service
[Unit]
Description=sync google calendar

AssertPathExists=/home/%u/wiki
AssertPathExists=/home/%u/priv

After=mount-priv.service

# gcal accesses and changes ~/wiki
Requisite=synced.target
After=synced.target

Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
Slice=workers.slice

SyslogIdentifier=gcal-sync
ExecStart=/usr/bin/make -C "${AB2_WIKI_DIR}" ical
ExecStart=/usr/bin/gcal -S
