# user/gcont-sync.service
[Unit]
Description=sync google contacts

Requires=mount-priv.service
After=mount-priv.service

# gcont store contacts in ~/mail
Requisite=initial-synced-all.service
After=initial-synced-all.service

Requires=network-online.service
After=network-online.service

[Service]
Type=oneshot
ExecStart=/usr/bin/gcont -d
