# user/git-gc@.service
[Unit]
Description=garbage collect git repo at %I

ConditionPathExists=/home/%u/%I/.git

[Service]
Type=oneshot

SyslogIdentifier=git-gc-%i
ExecStart=/usr/bin/git -C "/home/%u/%I" gc --auto 6700
