# user/glances.service
[Unit]
Description=glances monitor dropdown

PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
Slice=apps.slice

ExecStart=/usr/bin/termite --title=glances -e "/usr/bin/glances"

