# user/gmpc.service
[Unit]
Description=mpd client

Requisite=desktop.target
After=desktop.target

[Service]
Type=simple
ExecStart=/usr/bin/gmpc --debug-level=1
