# user/gnoti.service
[Unit]
Description=Gmail notifier

[Service]
Type=simple
ExecStart=/usr/bin/gnoti

[Install]
WantedBy=work.target
