# user/goldendict.service
[Unit]
Description=Goldendict

Requisite=wm.target
After=wm.target

[Service]
Type=simple
SyslogIdentifier=goldendict
Slice=apps.slice

ExecStart=/usr/bin/goldendict

