# user/gpg-agent.service
[Unit]
Description=gpg-agent daemon

ConditionHost=!skynet
ConditionHost=!quark
ConditionPathExists=/home/%u/.gnupg

Before=daemons.target

[Service]
Type=forking
Slice=daemons.slice
ExecStart=/usr/bin/gpg-agent --daemon

# give it some time
ExecStartPost=/usr/bin/sleep 1

[Install]
WantedBy=daemons.target