# user/gpg-agent-forward.socket
[Unit]
Description=forward ssh agent requests

[Socket]
ListenStream=%h/.gnupg/S.gpg-agent.forward

[Install]
WantedBy=sockets.target