# user/gpg-agent-forward.socket
[Unit]
Description=listen to gpg-agent requests from the net

[Socket]
ListenStream=11000

[Install]
WantedBy=sockets.target