# user/gpg-agent-auth.service
[Unit]
Description=refresh gpg-agent passphrase

Requires=gpg-agent.service
After=gpg-agent.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo true

