# user/gpg-agent.service
[Unit]
Description=gpg-agent daemon

[Service]
Type=forking
ExecStart=/usr/bin/gpg-agent --daemon --use-standard-socket
