# user/gpg-agent-forward.socket
[Unit]
Description=forward gpg agent socket

[Socket]
ListenStream=%t/gpg-agent/S.gpg-agent.forward
Service=gpg-agent-forward.service
SocketMode=755

[Install]
WantedBy=sockets.target