# user/gpg-agent-forward@.service
[Unit]
Description=forward gpg-agent to %i

Requires=gpg-agent.service
After=gpg-agent.service

[Service]
Type=simple
SyslogIdentifier=gpg-agent-proxy
ExecStart=/usr/bin/ssh -o ControlMaster=no -o ControlPath=none -R /home/abdo/.gnupg/S.gpg-agent:/home/abdo/.gnupg/S.gpg-agent -R /home/abdo/.gnupg/S.gpg-agent.ssh:/home/abdo/.gnupg/S.gpg-agent.ssh %i -N

