# user/gpg-agent-forward.socket
[Unit]
Description=forward gpg agent socket

[Socket]
ListenStream=%t/gpg-agent/S.gpg-agent.forward

[Install]
WantedBy=sockets.target