# user/gpg-agent.socket
[Unit]
Description=gpg agent socket

[Socket]
ListenStream=%t/gpg-agent/S.gpg-agent
Service=gpg-agent-proxy.service

[Install]
WantedBy=sockets.target