# user/gpg-agent-auth.service
[Unit]
Description=refresh gpg-agent passphrase
Requisite=gpg-agent.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo true

