# user/gpg-agent.service
[Unit]
Description=gpg-agent daemon

After=xorg.target
Requisite=xorg.target
StopWhenUnneeded=true

[Service]
Type=forking
ExecStart=/usr/bin/gpg-agent --daemon --use-standard-socket
