# user/gpg-agent-forward.service
[Unit]
Description=forward gpg agent requests over the net

[Service]
Type=simple
SyslogIdentifier=gpg-agent-forward
ExecStart=/usr/lib/systemd/systemd-socket-proxyd localhost:11010

