# user/gpg-agent-pause.service
[Unit]
Description=pause local gpg-agent
StopWhenUnneeded=true

[Service]
Type=oneshot
RemainAfterExit=yes

StandardOutput=null
ExecStart=/usr/bin/systemctl --user start gpg-agent-remote.service
ExecStop=/usr/bin/systemctl --user start gpg-agent.service gpg-agent-forward.service
