# user/gpg-agent.service
[Unit]
Description=gpg-agent daemon
ConditionHost=!skynet


[Service]
Type=forking
ExecStart=/home/abdo/.config/systemd/scripts/gpg-agent.sh start
