# user/gpg-agent-ssh-forward.socket
[Unit]
Description=forward ssh agent requests

[Socket]
ListenStream=%t/gpg-agent/S.gpg-agent.ssh.forward

[Install]
WantedBy=sockets.target