# user/gpg-agent.service
[Unit]
Description=GPG Agent

ConditionPathIsDirectory=/home/%u/.gnupg

[Service]
Type=forking
ExecStart=/usr/bin/gpg-agent --daemon --homedir=%h/.gnupg 
Restart=on-failure

[Install]
WantedBy=default.target