# user/gpg-agent.target
[Unit]
Description=gpg-agent daemon plus redirections
BindsTo=gpg-agent.service
Wants=gpg-agent-forward.socket gpg-agent-ssh-forward.socket
Wants=gpg-agent-proxy.service gpg-agent-ssh-proxy.service