# user/gpg-agent-ssh.socket
[Unit]
Description=ssh agent socket

[Socket]
ListenStream=%t/gpg-agent/S.gpg-agent.ssh
Service=gpg-agent-ssh-proxy.service
SocketMode=755

[Install]
WantedBy=sockets.target