# user/graphical.target
[Unit]
Description=graphical environment

RefuseManualStart=true
StopWhenUnneeded=true

Conflicts=console.target
After=console.target
