# user/graphical.target
[Unit]
Description=graphical server

RefuseManualStart=yes

Requires=user.target
After=user.target

After=xorg.service vncserver.service


