# user/graphical-daemons.target
[Unit]
Description=graphical daemons

Requisite=graphical.target
After=graphical.target
