# user/graphical.target
[Unit]
Description=graphical environment

After=xserver.target

RefuseManualStart=true
StopWhenUnneeded=true

Conflicts=console.target
