# user/graphical-daemons.target
[Unit]
Description=graphical daemons

Requires=graphical.target
After=graphical.target

