# user/graphical.target
[Unit]
Description=graphical server

RefuseManualStart=true
StopWhenUnneeded=true

