# user/graphical.target
[Unit]
Description=graphical server

RefuseManualStart=true
StopWhenUnneeded=true

Requires=user.target
After=user.target



