# user/hi.target
[Unit]
Description=trigger initial actions

After=initial-synced-data.target initial-synced-mail.target initial-synced-system.target
Requires=initial-synced-data.target initial-synced-mail.target
