# user/hi.target
[Unit]
Description=trigger initial actions
StopWhenUnneeded=true

OnFailure=journal-mailer.service

Requires=async-default.target
After=async-default.target

Requires=fetch.target
After=fetch.target