# user/hi.target
[Unit]
Description=trigger initial actions
StopWhenUnneeded=true
