# user/hi.target
[Unit]
Description=initial synchronization

StopWhenUnneeded=true
AllowIsolate=true
