# user/hodge-suspend.timer
[Unit]
Description=suspend hodge at night
ConditionHost=hodge

[Timer]
OnCalendar=20:00
Unit=sync-byez.service

[Install]
WantedBy=timers.target