# user/homedir-sync@.service
[Unit]
Description=sync homedir to remote %i

ConditionPathExists=/home/%u/.mrconfig

Requires=mount-priv.service
After=mount-priv.service

[Service]
Type=oneshot
Slice=workers.slice
SyslogIdentifier=homedir-sync

WorkingDirectory=/home/%u
ExecStart=/usr/bin/systemctl --user start presync.target
ExecStart=/usr/bin/mr --stats -t sync %i
ExecStart=/usr/bin/systemctl --user start synced.target

[Install]
Alias=sync@.service
