# user/homedir-sync@.service
[Unit]
Description=sync homedir to remote %i

AssertPathExists=/home/%u/.mrconfig
AssertPathExists=/home/%u/priv

After=mount-priv.service

[Service]
Type=oneshot
Slice=workers.slice
SyslogIdentifier=homedir-sync

WorkingDirectory=/home/%u
ExecStart=/home/abdo/etc/mr/bin/homedir sync %i

[Install]
Alias=sync@.service
