# user/host-wait-online@.service
[Unit]
Description=wait until host %i is online

Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
SyslogIdentifier=host-wait
Slice=daemon.slice

ExecStart=/usr/bin/sh -c "while true; do ping -q -c 1 %i > /dev/null && break; sleep 10; done" 
