# user/housekeeping.target
[Unit]
Description=houskeeping timers

ConditionACPower=yes

[Install]
WantedBy=udev-ac.target