# user/impulse.service
[Unit]
Description=novation impulse power
ConditionHost=grothendieck

[Service]
Type=oneshot
SyslogIdentifier=impulse
RemainAfterExit=true

ExecStart=/home/abdo/.config/systemd/scripts/usbpower.sh "on" "6-2"
ExecStop=/home/abdo/.config/systemd/scripts/usbpower.sh "off" "6-2"


