# user/impulse.service
[Unit]
Description=novation impulse keyboard

# Alias set in udev rules
BindsTo=sys-subsystem-sound-devices-impulse.device

[Service]
Type=oneshot

SyslogIdentifier=impulse
RemainAfterExit=true

ExecStart=/usr/bin/sudo systemctl start impulse.service
ExecStop=/usr/bin/sudo systemctl stop impulse.service


