# user/initial-sync.target
[Unit]
Description=initial synchronization

StopWhenUnneeded=true
AllowIsolate=true

Requires=data-synced.target fetch.target
After=data-synced.target fetch.target
