# user/initial-synced-data.target
[Unit]
Description=initial data sync

RequiresOverridable=async@default.service
After=async@default.service

