# user/initial-synced-fake.service
[Unit]
Description=fake synced state

Requires=mount-priv.service
After=mount-priv.service

[Service]
Type=oneshot
SyslogIdentifier=async-fake
ExecStart=/usr/bin/systemctl --user --job-mode=ignore-requirements start initial-synced-all.target
