# user/initial-synced-system.target
[Unit]
Description=initial system sync

Requires=pacman-update.service
After=pacman-update.service

