# user/initial-synced-base.target
[Unit]
Description=initial base sync

Requires=async-default-base.service
After=async-default-base.service

