# user/initial-synced-mail.target
[Unit]
Description=initial mail sync

AllowIsolate=true

Requires=initial-synced-data.target

BindsTo=fetch.target
After=fetch.target

