# user/initial-synced-work.target
[Unit]
Description=initial work sync

Requires=async-default-work.service
After=async-default-work.service

