# user/initial-synced-rest.target
[Unit]
Description=initial data sync

Requires=async-default-rest.service
After=async-default-rest.service

