# user/initial-synced-system.target
[Unit]
Description=initial system sync

AllowIsolate=true

BindsTo=pacman-update.service
After=pacman-update.service

