# user/initial-synced-lib.target
[Unit]
Description=initial work sync

Requires=async-default-lib.service
After=async-default-lib.service

