# user/initial-synced-mail.target
[Unit]
Description=initial mail sync

Requires=initial-synced-data.target
After=initial-synced-data.target

BindsTo=fetch.target
After=fetch.target

