# user/jack.service
[Unit]
Description=jack server
Documentation=man:jack(1)

Requires=jackdbus.service
After=jackdbus.service

PartOf=jack.target
Before=jack.target

Requires=audio.slice
After=audio.slice

[Service]
Type=oneshot
RemainAfterExit=true

# we use jackdbus
ExecStart=/usr/bin/jack_control start
ExecStop=/usr/bin/jack_control stop

[Install]
WantedBy=jack.target