# user/jack.target
[Unit]
Description=jack server

BindsTo=jack.service
After=jack.service
