# user/jack.target
[Unit]
Description=jack sound server

Requires=pulseaudio.service
After=pulseaudio.service



