# user/jack.target
[Unit]
Description=Jack server

# We let pulse handle everything and then jack steals the main device.
Requires=pulseaudio.service