# user/journal-triggerd.service
[Unit]
Description=triggers on journal messages

[Service]
Type=simple
ExecStart=/usr/bin/journal-triggerd "%h/.config/journal-triggerd"

