# user/journal-dropdown.service
[Unit]
Description=journal syslog terminal

[Service]
Type=simple
SyslogIdentifier=syslog
ExecStart=/usr/bin/termite -e "jctl -f -n50"

Slice=dropdown.slice