# user/jwatch.service
[Unit]
Description=journal monitor

[Service]
Type=simple
ExecStart=/usr/bin/jwatch trigger -f

