# user/jwatch.service
[Unit]
Description=journal monitor

[Service]
Type=simple
Slice=daemons.slice
ExecStart=/usr/bin/jwatch trigger -f
