# user/landscape.target
[Unit]
Description=landscape orientation
Conflicts=portrait.target

Requires=yoga-rotate@0.service