# user/laptop.target
[Unit]
Description=laptop configuration
AllowIsolate=true

Conflicts=console.target desktop.target tablet.target vncserver.target

