# user/linuxsampler-conf@.service
[Unit]
Description=linuxsampler config %i

Requisite=jack.service
After=jack.service

Requires=audio.slice
After=audio.slice

Requires=linuxsampler.service
After=linuxsampler.service


[Service]
Type=oneshot

SyslogIdentifier=linuxsampler-%i
ExecStart=/usr/bin/bash -c "cat ${AB2_CONF_DIR}/audio/linuxsampler/%i.lscp | nc localhost 8888 || true"

Slice=audio.slice

