# user/lock.target
[Unit]
Description=lock screen

BindsTo=i3lock.service
After=i3lock.service

Conflicts=mount-priv.service