# user/lock.target

[Unit]
Description=Locked screen
Requires=dbus.socket
BindsTo=slimlock.service
AllowIsolate=true
