# user/lock.target
[Unit]
Description=lock screen

# bind to slimlock, so that when unlocking lock.target goes down
BindsTo=slimlock.service
After=slimlock.service
