# user/lock.target

[Unit]
Description=Locked screen
Requires=dbus.socket
BindsTo=slimlock@0.service
AllowIsolate=true
