# user/lock.target

[Unit]
Description=Locked screen
Requires=dbus.socket
Requires=xorg.service
BindsTo=slimlock.service
AllowIsolate=true
