# user/lock.target

[Unit]
Description=Locked screen
Requires=dbus.socket
AllowIsolate=true
