# user/machine-mount@.service
[Unit]
Description=mount disk image ~/virt/qemu/%i.qcow2 vm

Conflicts=machine-%i.service
After=machine-%i.service

Requires=machine-nbd@%i.service
After=machine-nbd@%i.service

[Service]
Type=oneshot
RemainAfterExit=true
Slice=machine.slice

ExecStartPre=/usr/bin/sudo mkdir -p /media/%i
ExecStart=/usr/bin/sudo mount /dev/nbd0p1 /media/%i

ExecStop=/usr/bin/sudo umount /media/%i
