# user/machine-vnc@.service
[Unit]
Description=vnc client for %i vm

[Service]
Type=simple
Slice=machines.slice

EnvironmentFile=/home/abdo/virt/etc/machines.conf
EnvironmentFile=/home/abdo/virt/etc/machines.d/%i.conf

ExecStart=/usr/bin/vncviewer localhost::$PORT_VNC
