# user/machine-mount@.service
[Unit]
Description=mount disk of %i vm

Conflicts=machine-%i.service

[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=machine-mount-%i

ExecStart=/usr/bin/sudo qemu-nbd -c /dev/nbd0 ${AB2_VIRT_DIR}/qemu/%i.qcow2
ExecStart=/usr/bin/mkdir -p /media/%i
ExecStart=/usr/bin/mount /dev/nbd0p1 /media/%i

ExecStop=/usr/bin/umount /media/%i
ExecStop=/usr/bin/sudo qemu-nbd -d /dev/nbd0
