# user/machines.slice
[Unit]
Description=virtual machines

Before=slices.target
