# user/mail-db.service
[Unit]
Description=mail database

StopWhenUnneeded=true
ConditionPathExists=${AB2_MAIL_DIR}

Requisite=initial-synced-all.service
After=initial-synced-all.service

[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=mail-db
ExecStart=/usr/bin/true
ExecStop=/usr/bin/systemctl --user start commit-mail.service


