# user/mail-db.service
[Unit]
Description=mail database

StopWhenUnneeded=true
Conflicts=mail-commit.service

[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=mail-db
ExecStart=/usr/bin/true
ExecStop=/home/abdo/.config/systemd/scripts/mu4e.sh commit

