# user/mail-commit.service
[Unit]
Description=commit mail

Requires=mail-stop.service

[Service]
Type=oneshot

SyslogIdentifier=mail-commit
ExecStart=/home/abdo/.config/systemd/scripts/mu4e.sh commit
