# user/mail.target
[Unit]
Description=mail and news timers
Requires=dbus.socket
StopWhenUnneeded=true
AllowIsolate=true
