# user/mail-db.service
[Unit]
Description=mail database
StopWhenUnneeded=true

Requisite=initial-synced-all.target
After=initial-synced-all.target

[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=mail-db
ExecStart=/usr/bin/true
ExecStop=/usr/bin/systemctl --user start commit-mail.service


