# user/mail.timer
[Timer]
OnUnitInactiveSec=3min
OnActiveSec=3min
Unit=mail.service
