# user/mount-priv.service
[Unit]
Description=mount priv

ConditionPathExists=${HOME}/.ecryptfs/priv.conf

Wants=gpg-agent.service
After=gpg-agent.service

[Service]
Type=oneshot
RemainAfterExit=true

SyslogIdentifier=priv
ExecStart=/usr/bin/bash -c "gpg --batch -d ${HOME}/.ecryptfs/priv-passphrase 2> /dev/null | ecryptfs-add-passphrase --fnek"
ExecStart=/usr/bin/mount.ecryptfs_private priv

ExecStop=/usr/bin/umount.ecryptfs_private priv
