# user/mpc-pause.service
[Unit]
Description=Run mpc pause
Requisite=mpd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mpc pause
