# user/mpc-pause.service
[Unit]
Description=run mpc pause
Requisite=mpd.service
StopWhenUnneeded=true

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=/usr/bin/mpc pause
