# user/mpc@.service
[Unit]
Description=Run mpc command
After=mpd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mpc %I 
