# user/mpd.socket
[Unit]
Description=Music Player Daemon socket

# NOTE: Do not add After dependencies here to avoid a dependency cycle on shutdown:
# synced.target > basic.target > sockets.target > mpd.socket > synced.target
Requisite=synced.target

[Socket]
ListenStream=%t/mpd/mpd.socket

# gmpc doesn't know about unix sockets...
ListenStream=6600

[Install]
WantedBy=synced.target
