# user/mpd.service
[Unit]
Description=Music Player Daemon
StopWhenUnneeded=true
After=pulseaudio.service

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p /run/user/%U/mpd
# ExecStart=/usr/bin/mpd
# PIDFile=/run/user/%U/mpd/pid
ExecStart=/usr/bin/mpd --no-daemon

[Install]
WantedBy=audio.target
