# user/mpd.service
[Unit]
Description=music player daemon

Requires=mpd.socket
Wants=mpdscribble.service

PartOf=audio.target

Requires=pulseaudio.service
After=pulseaudio.service

[Service]
Type=simple
Environment=PULSE_PROP=application.icon_name=emblem-sound
ExecStart=/usr/bin/mpd --no-daemon

[Install]
WantedBy=audio.target
