# user/mpdscribble.service
[Unit]
Description = MPD Scribbler

[Service]
Type = simple
ExecStart = /usr/bin/mpdscribble -D
Restart = always

[Install]
WantedBy = audio.target
