# user/mpdscribble.service
[Unit]
Description=mpd scrobbbler

Requisite=synced.target
After=synced.target

Requires=mpd.socket mpd.service
After=mpd.socket mpd.service


[Service]
Type = simple
Slice=daemons.slice

ExecStart = /usr/bin/mpdscribble -D --conf %h/etc/mpd/mpdscribble.conf

[Install]
WantedBy=mpd.service
