# user/mu-server.service
[Unit]
Description=mu server

Requires=mu-server.socket
After=mu-server.socket

[Service]
Type=simple

Slice=daemons.slice

StandardInput=socket
ExecStart=/usr/bin/mu server


