# user/mu-server.socket
[Unit]
Description=socket for mu server

[Socket]
ListenStream=%t/mu/mu.socket

[Install]
WantedBy=sockets.target
