# user/mu4e.service
[Unit]
Description=mu4e mail client

Requisite=desktop.target
After=desktop.target

PartOf=mail.target

[Service]
Type=simple
SyslogIdentifier=mu4e
Slice=dropdown.slice

ExecStart=/usr/bin/emacs -mail --title "emacs-mail"
KillSignal=SIGINT

