# user/music-synth.target
[Unit]
Description=play synthetizer

Requires=jack.target
After=jack.target

Wants=fluidsynth.service qjackctl-play.service
After=fluidsynth.service qjackctl-play.service
