# user/music-piano.target
[Unit]
Description=play synth piano

Requires=jack.service
After=jack.service

Wants=pianoteq.service
After=pianoteq.service
