# user/music-pianoteq.target
[Unit]
Description=play pianoteq

Requires=jack.service
After=jack.service

Wants=pianoteq.service
After=pianoteq.service
