# user/music-piano.target
[Unit]
Description=play sampled piano

Requires=jack.target
After=jack.target

Wants=linuxsampler.service linuxsampler-conf@piano.service qjackctl-play.service
After=linuxsampler.service linuxsampler-conf@piano.service qjackctl-play.service
