# user/music-synth.target
[Unit]
Description=play synthetizer

Requires=jack.service
After=jack.service

Wants=fluidsynth.service qjackctl-play.service
After=fluidsynth.service qjackctl-play.service
