# user/music-studio.target
[Unit]
Description=music midi studio

Wants=jack.target impulse.service
Wants=qjackctl-studio.service qtractor.service yoshimi.service fluidsynth.service
