# user/mutag-autotag.service
[Unit]
Description=postprocess mail and news

Requisite=initial-synced-all.target
After=initial-synced-all.target

After=offlineimap.service
After=rss2maildir.service
Conflicts=mu4e.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mutag --index --silent
ExecStart=/usr/bin/mutag --autotag --modified --update --silent -q 'date:1m..now'
ExecStart=/usr/bin/mutag --index --silent
