# user/mutag-autotag.service
[Unit]
Description=postprocess mail and news

ConditionPathExists=/home/%u/mail
ConditionPathExists=/usr/bin/mutag

Requisite=synced.target
After=synced.target

After=offlineimap.service rss2maildir.service
Conflicts=mu4e.service

[Service]
Type=oneshot
SyslogIdentifier=mutag-autotag
ExecStart=/usr/bin/mutag --index --silent
ExecStart=/usr/bin/mutag --autotag --modified --update --silent -q 'date:1m..now'

