# user/mutag-queue.path
[Unit]
Description=watch mail queue
ConditionHost=!skynet

[Path]
PathChanged=/home/abdo/.msmtp.queue
Unit=mutag-newmail.service
