# user/mutag-mail.path
[Unit]
Description=watch mail database changes

ConditionPathExists=${AB2_MAIL_DIR}
ConditionPathExists=${HOME}/.mu/xapian

[Path]
PathChanged=/home/abdo/.mu/xapian
Unit=mutag-newmail.service
