# user/mutag-mail.path
[Unit]
Description=watch mail database changes
ConditionHost=!skynet

[Path]
PathChanged=/home/abdo/.mu/xapian
Unit=mutag-newmail.service
