# user/mutag-queue.path
[Unit]
Description=watch mail queue

ConditionPathExists=/home/abdo/mail
ConditionPathExists=/usr/bin/mutag

[Path]
PathChanged=/home/abdo/mail/queue
Unit=mutag-newmail.service
