# user/mutag-queue.path
[Unit]
Description=watch mail queue

ConditionPathExists=${AB2_MAIL_DIR}

[Path]
PathChanged=/home/abdo/mail/queue
Unit=mutag-newmail.service
