# user/mutag-index.service
[Unit]
Description=index mail and news

Requisite=initial-synced-all.service
After=initial-synced-all.service

After=offlineimap.service
After=rss2maildir.service
Conflicts=mu4e.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mutag --index --silent
