# user/ncmpcpp.service
[Unit]
Description=ncmpcpp mpd client

PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
Slice=apps.slice

ExecStart=/usr/bin/termite --title=musicplayer -e "/usr/bin/ncmpcpp -h %t/mpd/mpd.socket"

[Install]
Alias=musicplayer.service
