# user/network-online.target
[Unit]
Description=network is online
DefaultDependencies=no
Conflicts=suspend.target poweroff.target reboot.target
