# user/news.timer
[Timer]
OnUnitInactiveSec=2h
OnActiveSec=2h
Unit=news.service
