# user/noop.service
[Unit]
Description=noop

[Service]
Type=simple
ExecStart=/usr/bin/true
