# user/notes-dropdown.service
[Unit]
Description=notes dropdown

Requisite=wm.target
After=wm.target

Requires=git-use@work-wiki.service
After=git-use@work-wiki.service

[Service]
Type=simple
SyslogIdentifier=notes
Slice=dropdown.slice

ExecStart=/usr/bin/emacs -notes --title "emacs-notes"
KillSignal=SIGINT

