# user/notes.service
[Unit]
Description=chat

Requisite=wm.target
After=wm.target

Requires=orgmode-db.service
After=orgmode-db.service

[Service]
Type=simple
SyslogIdentifier=notes
Slice=dropdown.slice

ExecStart=/usr/bin/emacs -notes --title "emacs-notes"
KillSignal=SIGINT

