# user/offlineimap.service
[Unit]
Description=offlineimap

ConditionPathExists=/home/%u/mail
ConditionPathExists=/home/%u/.netrc

Requires=mount-priv.service
After=mount-priv.service

Requisite=synced.target
After=synced.target

Requires=network-online.service
After=network-online.service

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap
