# user/offlineimap.service
[Unit]
Description=offlineimap

Requisite=initial-synced-data.target
After=initial-synced-data.target
After=journal-mailer.service

Requires=online.service
After=online.service

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap -l %h/var/log/offlineimap.log
