# user/offlineimap.service
[Unit]
Description=offlineimap

Requires=mount-priv.service
After=mount-priv.service

Requisite=initial-synced-all.target
After=initial-synced-all.target

Requires=online.service
After=online.service

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap
