# user/offlineimap.service
[Unit]
Description=offlineimap

AssertPathExists=/home/%u/mail
AssertPathExists=/home/%u/.netrc

Requisite=synced.target
After=synced.target

Requires=git-use@mail.service
After=git-use@mail.service

Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
Slice=workers.slice

ExecStart=/usr/bin/offlineimap
