# user/offlineimap.service
[Unit]
Description=offlineimap

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap -l %h/var/log/offlineimap.log
