# user/openvpn@.service
[Unit]
Description=OpenVPN connection to %I

After=machine@%i.service

[Service]

Type=simple

# start system service and wait for it
ExecStartPre=/usr/bin/sudo systemctl --system start openvpn@%i.service

# synchronously wait until that service stops
ExecStart=/usr/bin/sudo systemctl --system start --wait openvpn@%i.service

ExecStop=/usr/bin/sudo systemctl --system stop openvpn@%i.service
