# user/openbox.target
[Unit]
Description=openbox desktop environment
Conflicts=console.target

BindsTo=openbox.service
After=openbox.service
