# user/openvpn@.service
[Unit]
Description=OpenVPN connection to %I

After=machine@%i.service

[Service]

Type=simple

# Start openvpn and synchronously wait until it stops
ExecStart=/usr/bin/sudo systemctl --system start --wait openvpn-client@%i.service

ExecStop=/usr/bin/sudo systemctl --system stop openvpn-client@%i.service
