# user/openbox.target

[Unit]
Description=A Desktop Environment built on Openbox
Requires=dbus.socket
AllowIsolate=true
