# user/openbox@.service
[Unit]
Description=Openbox window manager at display %i

After=xserver@%i.target
After=xorg@%i.service vncserver@%i.service x11vnc@%i.service xephyr@%i.service
After=compton.service xrandr-auto.service

Before=tint2@%i.service
Before=desktop.target

Requisite=xserver@%i.target
Requires=tint2@%i.service
PartOf=wm@openbox.target

[Service]
PAMName=login
TTYPath=/dev/tty7
Type=simple
SyslogIdentifier=openbox
Environment=DISPLAY=:%i
ExecStart=/home/abdo/.config/systemd/scripts/openbox.sh start
