# user/openbox@.service
[Unit]
Description=openbox window manager at display %i

After=xserver@%i.target
After=xorg@%i.service vncserver@%i.service x11vnc@%i.service xephyr@%i.service
After=compton.service xrandr-auto.service

Before=tint2@%i.service
Before=desktop.target

Requisite=xserver@%i.target
Requires=tint2@%i.service
PartOf=wm@openbox.target

[Service]
Type=simple
SyslogIdentifier=openbox
Slice=wm.slice

Environment=DISPLAY=:%i
ExecStart=/usr/bin/openbox
KillSignal=SIGINT
TimeoutSec=15
