# user/openbox@.service
[Unit]
Description=Openbox Window Manager
After=xorg.target
StopWhenUnneeded=true

Requires=xorg.target
Requires=dbus.socket

[Service]
EnvironmentFile=%h/.config/systemd/conf/xenvironment
Environment=DISPLAY=:%I
ExecStart=/home/abdo/.config/systemd/scripts/openbox start
# ExecStop=/home/abdo/.config/systemd/scripts/openbox stop

[Install]
WantedBy=openbox.target
