# user/openbox.service
[Unit]
Description=Openbox Window Manager
After=xorg.target
StopWhenUnneeded=true
Requires=xorg.service
Requires=dbus.socket
Conflicts=awesome.service xmonad.service

[Service]
EnvironmentFile=%h/.config/systemd/conf/xenvironment
ExecStart=/home/abdo/.config/systemd/scripts/openbox start
# ExecStop=/home/abdo/.config/systemd/scripts/openbox stop

[Install]
WantedBy=openbox.target
