# user/orgmode.service
[Unit]
Description=org-mode

Requisite=graphical.target
After=graphical.target

Requires=git-use@wiki.service
After=git-use@wiki.service

[Service]
Type=simple
Slice=dropdown.slice

ExecStart=/usr/bin/emacs -org --title "emacs-org"
KillSignal=SIGINT


