# user/orgmode-save.service
[Unit]
Description=save orgmode buffers

# So that orgmode.service is stoped after this. The stopping order gets reversed.
After=orgmode.service

[Service]
Type=oneshot

# NOTE: If orgmode.service is not active, always succeed since there is nothing to be done.
#       It would be nice to write the dependency on orgmode.service directly in systemd!
ExecStart=/home/abdo/bin/emacs-run -q -s orgmode.service -k orgmode '(save-some-buffers t)'

