# user/orgmode.service
[Unit]
Description=org-mode
After=xorg.target

Requisite=xorg.target

[Service]
Type=simple
SyslogIdentifier=orgmode

ExecStart=/home/abdo/.config/systemd/scripts/orgmode.sh start
ExecStop=/home/abdo/.config/systemd/scripts/orgmode.sh stop

# commit changes
ExecStopPost=/home/abdo/.config/systemd/scripts/orgmode.sh commit

