# user/patchage.service
[Unit]
Description=jack patch bay

Requisite=jack.service jackdbus.service

[Service]
Type=simple
SyslogIdentifier=patchage
Slice=apps.slice

ExecStart=/usr/bin/patchage


