# user/paur-update.service
[Unit]
Description=update packages
Requires=online.service
After=online.service

[Service]
Type=oneshot
ExecStart=/usr/bin/paur -Syu
