# user/pianoteq-gui.service
[Unit]
Description=pianotec synth GUI

Conflicts=pianoteq.service

Requisite=wm.target
After=wm.target

Requires=jack.service
After=jack.service

Requires=audio.slice
After=audio.slice

[Service]
Type=simple
ExecStart=/usr/bin/pianoteq

# realtime scheduling
CPUSchedulingPolicy=rr
CPUSchedulingPriority=60

Slice=audio.slice
