# user/plank.service
[Unit]
Description=Plank dock

Requisite=graphical-session.target
After=graphical-session.target

PartOf=tablet.target

[Service]
Type=simple

StandardOutput=null
ExecStart=/usr/bin/plank
