# user/poweroff.target
[Unit]
Description=poweroff system
DefaultDependencies=no

Requires=systemd-poweroff.service

