# user/poweroff.target
[Unit]
Description=poweroff system
DefaultDependencies=no
BindsTo=systemd-poweroff.service
After=systemd-poweroff.service
