# user/power-ac.target
[Unit]
Description=AC Power

BindsTo=sys-subsystem-power_supply-devices-AC.device
