# user/presentation.target
[Unit]
Description=presentation configuration
AllowIsolate=true

Wants=laptop.target

Conflicts=console.target desktop.target
