# user/presync.target
[Unit]
Description=homedir presync hook

Conflicts=synced.target