# user/pulseaudio.service
[Unit]
Description=PulseAudio Sound System
Documentation=man:pulseaudio(1)

PartOf=audio.target

Requires=audio.slice
After=audio.slice

[Service]
Type=dbus
BusName=org.pulseaudio.Server
BusName=org.PulseAudio1

ExecStart=/usr/bin/pulseaudio
SuccessExitStatus=0 1
Restart=always
RestartSec=1

Slice=audio.slice

[Install]
WantedBy=audio.target
