# user/pulseaudio.service
[Unit]
Description=PulseAudio Sound System
Documentation=man:pulseaudio(1)
Before=sound.target
Requires=dbus.socket

[Service]
Type=forking
BusName=org.pulseaudio.Server
ExecStart=/usr/bin/pulseaudio --start --daemonize --fail
PIDFile=/run/user/%U/pulse/pid
Restart=always
RestartSec=1

[Install]
WantedBy=audio.target
