# user/pulseaudio.service
[Unit]
Description=PulseAudio Sound System
Documentation=man:pulseaudio(1)

Requires=pulseaudio.socket
After=pulseaudio.socket

Before=audio.target
PartOf=audio.target

Requires=audio.slice
After=audio.slice

[Service]
Type=dbus

BusName=org.pulseaudio.Server
BusName=org.PulseAudio1

# realtime scheduling
CPUSchedulingPolicy=rr
CPUSchedulingPriority=60

ExecStart=/usr/bin/pulseaudio
SuccessExitStatus=0 1

Slice=audio.slice

[Install]
Also=pulseaudio.socket
WantedBy=audio.target