# user/qtractor.service
[Unit]
Description=qtractor midi sequencer

Requisite=jack.service
After=jack.service jackdbus.service jack.target

After=fluidsynth.service linuxsampler.service

[Service]
Type=simple
SyslogIdentifier=qtractor
Slice=apps.slice

# realtime scheduling
CPUSchedulingPolicy=rr
CPUSchedulingPriority=60

ExecStart=/usr/bin/qtractor


