# user/quit-wm.target
[Unit]
Description=quit window manager

Wants=console.target
Before=console.target
