# user/quit-wm.target
[Unit]
Description=quit window manager

Requires=console.service
Before=console.service

AllowIsolate=true
