# user/quit-wm.target
[Unit]
Description=quit window manager

Requires=console.target
Before=console.target
