# user/reboot.target
[Unit]
Description=reboot system
DefaultDependencies=no

Requires=systemd-reboot.service

