# user/redshift.service
[Unit]
Description=Redshift

After=xorg.target
After=compton.service
StopWhenUnneeded=true

Requisite=xorg.target

[Service]
Type=simple
ExecStart=/usr/bin/redshift -l 41.383:2.183 -t 6500:5500 -g 1.0 -m randr
