# user/rosegarden.service
[Unit]
Description=rosegarden midi sequencer

Requisite=desktop.target
After=desktop.target

Requires=jack.service
After=jack.service

After=fluidsynth.service linuxsampler.service

[Service]
Type=simple
SyslogIdentifier=rosegarden
Slice=apps.slice

# realtime scheduling
CPUSchedulingPolicy=rr
CPUSchedulingPriority=60

ExecStart=/usr/bin/rosegarden


