# user/rosegarden.service
[Unit]
Description=rosegarden midi sequencer

Requisite=wm.target
After=wm.target

Requires=jack.service
After=jack.service

After=fluidsynth.service linuxsampler.service

[Service]
Type=simple
Slice=apps.slice

# realtime scheduling
CPUSchedulingPolicy=rr
CPUSchedulingPriority=60

ExecStart=/usr/bin/rosegarden


