# user/rss2maildir.service
[Unit]
Description=rss2maildir

ConditionPathExists=/home/%u/mail

Requisite=synced.target
After=synced.target

Requires=network-online.service
After=network-online.service

[Service]
Type=oneshot
ExecStart=/usr/bin/rss2maildir -v 
