# user/rss2maildir.service
[Unit]
Description=rss2maildir

[Service]
Type=oneshot
ExecStart=/usr/bin/rss2maildir -v -l %h/var/log/rss2maildir.log
