# user/rtags.service
[Unit]
Description=rtags server

Requires=rtags.socket
After=rtags.socket

[Service]
Type=simple

# the default socket is in ~/.rdm
ExecStart=/usr/bin/rdm -v --inactivity-timeout 300




