# user/sage.service
[Unit]
Description=sage server

[Service]
ExecStart=/home/abdo/bin/sage -n

