# user/setenv.service
[Unit]
Description=Set the environment

Before=default.target

[Service]
Type=oneshot

# set environment variables from /etc/zshenv and ~/.zshenv
ExecStart=/usr/bin/sh -c '                            \
    cat /etc/zshenv ~/.zshenv |                       \
    sed -n "s/^\\s*export\\s*\\([^=]*\\)=.*$/\\1/p" | \
    xargs systemctl --user import-environment         \
'

# symlink $XDG_RUNTIME_DIR to ~/.runtime
ExecStart=/usr/bin/ln -sfT "%t" "%h/.runtime"
