# user/setenv.service
[Unit]
Description=Set the environment

[Service]
Type=oneshot
ExecStart=/home/abdo/.config/systemd/scripts/setenv.sh
