# user/sleep.target

[Unit]
Description=sleep state
AllowIsolate=true

# bind to lock.target so that when unlocking sleep.target goes down
BindsTo=lock.target
Requires=sleep.service

