# user/sleep.target
[Unit]
Description=sleep
DefaultDependencies=no
RefuseManualStart=yes
StopWhenUnneeded=yes
