# user/slices.target
[Unit]
Description=user slices

[Install]
WantedBy=default.target
