# user/ssh-proxy@.socket
[Unit]
Description=ssh connection proxy socket for %i

[Socket]
ListenStream=%t/ssh/%i.socket

[Install]
WantedBy=sockets.target
