# user/ssh-wait-online@.service
[Unit]
Description=wait until ssh host %i is online

Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
SyslogIdentifier=ssh-wait
Slice=daemons.slice

ExecStart=/usr/bin/sh -c "while true; do ssh -v %i true 2>&1 | grep -q 'Connection established' && break; sleep 10; done"
