# user/ssh@.target
[Unit]
Description=ssh connection to %i

Requires=sshmux@%i.service gpgfwd@%i.service
After=sshmux@%i.service gpgfwd@%i.service
