# user/sshmux@.service
[Unit]
Description=ssh master connection for %i

AssertHost=!%i

Conflicts=disconnect.target

Requires=gpg-agent.service
After=gpg-agent.service

Requires=network-online.target
After=network-online.target

[Service]
Type=forking
Slice=daemons.slice

SyslogIdentifier=sshmux
ExecStartPre=/usr/bin/mkdir -p "%t/ssh"
ExecStart=/usr/bin/ssh -M -o "ControlPersist=yes" %i true


