# user/stream@.service
[Unit]
Description=streaming client %i

Requires=network-online.target
After=network-online.target

Requisite=audio.target
After=audio.target

[Service]
Type=simple
SyslogIdentifier=stream-%i
Slice=audio.slice

ExecStart=/home/abdo/.config/systemd/scripts/stream.sh %i


