# user/suspend.target
[Unit]
Description=suspend system
DefaultDependencies=no
BindsTo=systemd-suspend.service
Before=systemd-suspend.service
