# user/suspend.target
[Unit]
Description=suspend system
DefaultDependencies=no

Requires=systemd-suspend.service
After=systemd-suspend.service
