# user/sway-session.target
[Unit]
Description=sway wm session
AllowIsolate=true

BindsTo=graphical-session.target
Before=graphical-session.target

BindsTo=sway.service
After=sway.service
