# user/sway-proxy.service
[Unit]
Description=proxy service to sway window manager

Conflicts=awesome.service openbox.service xmonad.service tint2.service

Requires=graphical.target
Before=graphical.target

[Service]
Type=oneshot
RemainAfterExit=true
Slice=wm.slice

ExecStart=/usr/bin/true
