# user/synced.target
[Unit]
Description=synced target
AllowIsolate=true