# user/sync-fake.service
[Unit]
Description=fake initial sync

Requires=mount-priv.service
After=mount-priv.service

[Service]
Type=oneshot
SyslogIdentifier=sync-fake
ExecStart=/usr/bin/systemctl --user --no-block --job-mode=ignore-requirements start sync-hi.service
