# user/sync-skynet.target
[Unit]
Description=synchronize with skynet

Requires=online.service
After=online.service

BindsTo=async@skynet.service
Before=async@skynet.service
