# user/sync-hypathia.service
[Unit]
Description=synchronize with hypathia sdcard

BindsTo=media-share-hypathia.mount
After=media-share-hypathia.mount

Requires=calibre-db.service

[Service]
Type=oneshot
ExecStart=/usr/bin/async -b sync hypathia --terse
