# user/synced.target
[Unit]
Description=synced
RefuseManualStart=true

BindsTo=flag@synced.service
After=flag@synced.service