# user/sync-hi.service
[Unit]
Description=initial sync

Requires=async-default-needed.service
After=async-default-needed.service

BindsTo=flag@synced.service
Before=flag@synced.service

OnFailure=journal-mailer.service


[Service]
Type=oneshot

SyslogIdentifier=sync-hi
ExecStart=/usr/bin/systemctl --user start synced.target                         
ExecStop=/usr/bin/systemctl --user stop synced.target                         
