# user/synced.target
[Unit]
Description=synced
AllowIsolate=true

BindsTo=flag@synced.service
After=flag@synced.service