# user/sync-hi.service
[Unit]
Description=initial homedir sync

AssertPathExists=/home/%u/arch/sync/Makefile
AssertPathExists=/home/%u/priv/mnt

Wants=journal-show.service
After=journal-show.service

OnFailure=journal-mailer.service

[Service]
Type=oneshot
Slice=workers.slice

Nice=19
WorkingDirectory=/home/%u/arch/sync
ExecStart=/usr/bin/make sync REMOTES=skynet,babel
ExecStart=/usr/bin/systemctl --user start synced.target
