# user/sync-byez.service
[Unit]
Description=final homedir sync and suspend

Conflicts=sync-hi.service sync-fake.service synced.target

Wants=journal-show.service
After=journal-show.service

Requires=async-default-all.service
After=async-default-all.service

OnFailure=journal-mailer.service


[Service]
Type=oneshot
SyslogIdentifier=sync-byez
ExecStart=/usr/bin/sleep 10
ExecStart=/usr/bin/systemctl --user start suspend.target
