# user/systemd-reboot.service
[Unit]
Description=reboot system

[Service]
RemainAfterExit=true
Type=oneshot

ExecStart=/usr/bin/sudo systemctl reboot

