# user/systemd-poweroff.service
[Unit]
Description=poweroff system

Requires=shutdown.target
After=shutdown.target


[Service]
RemainAfterExit=true
Type=oneshot

ExecStart=/usr/bin/sudo systemctl poweroff

