# user/systemd-poweroff.service
[Unit]
Description=poweroff system

[Service]
RemainAfterExit=true
Type=oneshot

ExecStart=/usr/bin/sudo systemctl poweroff

