# user/systemd-reboot.service

[Unit]
Description=reboot system

[Service]
Type=oneshot

ExecStart=/usr/bin/sudo systemctl reboot

