# user/systemd-delayed-poweroff.service
[Unit]
Description=poweroff system in 20 seconds

Requires=systemd-poweroff.service
Before=systemd-poweroff.service

[Service]
RemainAfterExit=true
Type=oneshot

ExecStart=/usr/bin/sleep 20

