# user/systemd-poweroff.service
[Unit]
Description=poweroff system

DefaultDependencies=no

[Service]
Type=oneshot

ExecStart=/usr/bin/sleep 1
# NOTE: don't want --force here like the system instance.
ExecStart=/usr/bin/sudo systemctl poweroff

