# user/systemd-suspend.service
[Unit]
Description=enter suspend state

BindsTo=lock.taraget
After=lock.target

Requires=sleep.target
After=sleep.target

[Service]
RemainAfterExit=true
Type=oneshot
ExecStart=/usr/bin/sudo systemctl suspend
