# user/systemd-poweroff.service

[Unit]
Description=poweroff system

[Service]
Type=oneshot

ExecStart=/usr/bin/sudo systemctl poweroff

