# user/sysupgrade.service
[Unit]
Description=upgrade system and reboot

Requires=pacman-update.service
After=pacman-update.service

OnFailure=journal-mailer.service

[Service]
Type=oneshot
SyslogIdentifier=sysupgrade
ExecStart=/usr/bin/systemctl --user start reboot.target