# user/tablet.target
[Unit]
Description=tablet configuration
AllowIsolate=true

Conflicts=console.target desktop.target laptop.target vncserver.target

