# user/tablet.target
[Unit]
Description=tablet configuration
AllowIsolate=true

# target conflicts
Conflicts=console.target desktop.target laptop.target

Requires=graphical.target
After=graphical.target

# service requirements
Requires=caribou.service plank.service touchegg.service
After=caribou.service plank.service touchegg.service