# user/termite@.service
[Unit]
Description=termite terminal %i

Requisite=desktop.target 
After=desktop.target

[Service]
Type=simple
SyslogIdentifier=termite
Slice=apps.slice

ExecStart=/usr/bin/termite


