# user/test.service
[Unit]
Description=a testing service

[Service]
Type=oneshot
SyslogIdentifier=test
ExecStart=/usr/bin/echo "Hello World"

[Install]
WantedBy=desktop.target