# user/test-B.service
[Unit]
Description=test B

[Service]
Type=simple
ExecStart=/usr/bin/echo BBB
