# user/test.timer
[Unit]
Description=test timer

[Timer]
OnUnitInactiveSec=30
Persistent=true
Unit=test.service
