# user/test-fail.service
[Unit]
Description=always fail

[Service]
Type=forking
ExecStart=/usr/bin/false
