# user/thunar-dropdown.service
[Unit]
Description=thunar dropdown

Requisite=wm.target 
After=wm.target

[Service]
Type=simple
Slice=dropdown.slice

ExecStart=/usr/bin/thunar


