# user/tint2.service
[Unit]
Description=tint2 panel

Requisite=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/tint2
