# user/tint2@.service
[Unit]
Description=Tint2 panel
After=xorg@%i.service vncserver@%i.service xephyr@%i.service

Requisite=xserver@%i.target
After=xserver@%i.target

[Service]
Type=simple
Environment=DISPLAY=:%i
ExecStart=/usr/bin/tint2
