# user/tmux.service
[Unit]
Description=tmux session
Documentation=man:tmux(1)

[Service]
Type=forking

SuccessExitStatus=0 1
RuntimeDirectory=tmux
ExecStart=/usr/bin/tmux -S "%t/tmux/default" new-session -d -s default
ExecStop=/usr/bin/tmux -S "%t/tmux/default" kill-session -t default

