# user/tmux-server.service
# Launches a tmux server with a default session.
# makes sure no previous tmux process exists 

[Unit]
Description=tmux default session
Documentation=man:tmux(1)

[Service]
Type=forking
EnvironmentFile=%h/.config/systemd/conf/xenvironment
ExecStartPre=/usr/bin/tmux kill-server
ExecStart=/usr/bin/tmux new-session -d -s default
ExecStop=/usr/bin/tmux kill-server
SucessExitStatus=0 1

[Install]
WantedBy=default.target

