# user/umount.target
[Unit]
Description=umount devices

DefaultDependencies=no
RefuseManualStart=yes

