# user/user.target
[Unit]
Description=user target
AllowIsolate=true

Wants=dbus.service
After=dbus.service