# user/vnc.target
[Unit]
Description=VNC Server
Requires=dbus.socket
BindsTo=vnc@2.service
StopWhenUnneeded=true
AllowIsolate=true
