# user/vnc@.target
[Unit]
Description=Window Manager
Requires=dbus.socket
BindsTo=vncserver@2.service
BindsTo=%i@2.service
AllowIsolate=true
