# user/vncserver.service
[Unit]
Description=vnc server

Conflicts=xorg.service weston.service weston.socket

Conflicts=console.target
After=console.target

BindsTo=xserver.target
Before=xserver.target

[Service]
Type=forking

ExecStart=/usr/bin/vncserver :0 -nolisten tcp -localhost -nevershared -geometry 800x600 -depth 24

[Install]
WantedBy=vncserver.target
