# user/vncserver.target
[Unit]
Description=VNC Server
Requires=dbus.socket
BindsTo=vncserver@2.service
StopWhenUnneeded=true
AllowIsolate=true
