# user/vncserver.target
[Unit]
Description=vncserver configuration
AllowIsolate=true

Requisite=graphical.target
After=graphical.target

Conflicts=console.target laptop.target tablet.target presentation.target desktop.target

