# user/wakeup.target
[Unit]
Description=wakeup from suspend
DefaultDependencies=no
