# user/wireguard@.service
[Unit]
Description=wireguard VPN for %i

[Service]
Type=oneshot
RemainAfterExit=true

# Start wireguard
ExecStart=/usr/bin/sudo systemctl --system start wireguard@%i.service

# Stop it
ExecStop=/usr/bin/sudo systemctl --system stop wireguard@%i.service
