# user/wm.slice
[Unit]
Description=window manager slice
