# user/wm.slice
[Unit]
Description=window manager

[Install]
WantedBy=wm.target
