# user/wm.slice
[Unit]
Description=window manager

Before=slices.target

[Install]
WantedBy=slices.target