# user/wm.slice
[Unit]
Description=window manager slice

Before=slices.target

[Install]
WanedBy=slices.target