# user/wm.target
[Unit]
Description=window manager
Conflicts=console.target

StopWhenUnneeded=true

After=awesome@0.service xmonad@0.service openbox@0.service