# user/wm.target
[Unit]
Description=window manager

RefuseManualStart=yes

Conflicts=console.target

Requisite=graphical.target
After=graphical.target
