# user/wm.target
[Unit]
Description=window manager
Conflicts=console.target quit-wm.target

Requires=graphical.target
After=graphical.target

After=awesome@0.service xmonad@0.service openbox@0.service