# user/wm@.target
[Unit]
Description=%i window manager via xorg

BindsTo=xorg.target
BindsTo=%i@0.service
AllowIsolate=true
