# user/wm@.target
[Unit]
Description=%i window manager via xorg

Conflicts=console.target
BindsTo=xorg.target
BindsTo=%i@0.service
After=xorg.target %i@0.service
