# user/wm@.target
[Unit]
Description=Window Manager
Requires=dbus.socket
BindsTo=xorg@0.service
BindsTo=%i@0.service
AllowIsolate=true
