# user/wol@.service
[Unit]
Description=wake on lan %i
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
SyslogIdentifier=wol
ExecStart=/usr/bin/wol -f /etc/ethers %i
