# user/work.target
[Unit]
Description=Work related servers
Requires=dbus.socket
StopWhenUnneeded=true
AllowIsolate=true
