# user/work.target
[Unit]
Description=Work related servers

StopWhenUnneeded=true
AllowIsolate=true
