# user/work.target
[Unit]
Description=Work related servers
After=desktop.target


StopWhenUnneeded=true
AllowIsolate=true
