# user/workers.slice
[Unit]
Description=user workers

Before=slices.target

[Install]
WantedBy=slices.target