# user/wpa-gui.service
[Unit]
Description=wpa_supplicant gui

ConditionHost=galois

Requisite=graphical.target
After=graphical.target

[Service]
Type=simple
Slice=wm.slice

ExecStart=/usr/bin/wpa_gui -t

[Install]
WantedBy=wm.target
