# user/x11vnc@.service
[Unit]
Description=X11VNC server at display %i
After=xorg@0.service

Requisite=xorg@0.service
BindsTo=xserver@%i.target
Before=xserver@%i.target
Conflicts=xorg@%i.service vncserver@%i.service xephyr@%i.service

[Service]
Type=forking
EnvironmentFile=%h/.config/systemd/conf/xenvironment
Environment=DISPLAY=:0
ExecStart=/usr/bin/x11vnc -usepw -forever -display :%i -rfbport 590%i -bg

