# user/x11vnc@.service
[Unit]
Description=X11VNC server at display %i
StopWhenUnneeded=true

BindsTo=xserver@%i.target
Before=xserver@%i.target
Conflicts=xorg@%i.service vncserver@%i.service xephyr@%i.service

[Service]
Type=forking
ExecStart=/usr/bin/x11vnc -usepw -forever -display :%i -rfbport 590%i -bg

