# user/x11vnc.target
[Unit]
Description=VNC Server attached to local X.org
Requires=dbus.socket
BindsTo=x11vnc-server@3.service
StopWhenUnneeded=true
AllowIsolate=true
