# user/x11vnc@.service
[Unit]
Description=X11VNC server at display %i

Requisite=xorg@%i.service
After=xorg@%i.service

Requires=xserver@%i.target
Before=xserver@%i.target

[Service]
Type=forking
ExecStart=/usr/bin/x11vnc -usepw -forever -display :%i -rfbport 5902 -bg


