# user/xephyr@.target
[Unit]
Description=Xephyr session

Requisite=xorg@0.service
Requires=dbus.socket
BindsTo=xephyr@3.service
BindsTo=%i@3.service
AllowIsolate=true
