# user/xephyr.service
[Unit]
Description=Xephyr server at display :1

Requires=xephyr.socket

PartOf=graphical-session.target
After=graphical-session.target


[Service]
Type=simple
ExecStart=/usr/bin/Xephyr -ac -br -noreset -screen 1024x768 :1

