# user/xinit.service
[Unit]
Description=X initialization commands
After=xorg.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/systemd/conf/xenvironment
ExecStart=/usr/bin/setxkbmap -option terminate:ctrl_alt_bksp 
ExecStart=/usr/bin/xset r rate 500 40  

[Install]
WantedBy=xorg.target
