# user/xmodmap@.service
[Unit]
Description=X modmap
Documentation=man:xmodmap(1)

Requisite=xserver@%i.target
After=xserver@%i.target

After=xorg@%i.service vncserver@%i.service xephyr@%i.service

[Service]
Type=oneshot
Environment=DISPLAY=:%i
ExecStart=/usr/bin/xmodmap %h/.config/systemd/conf/Xmodmap
