# user/xmodmap@.service
[Unit]
Description=X resource setting
Documentation=man:xmodmap(1)
After=xorg@%i.service
After=xinit@%i.service
BindsTo=xorg@%i.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/systemd/conf/xenvironment
Environment=DISPLAY=:%i
ExecStart=/usr/bin/xmodmap %h/.config/systemd/conf/Xmodmap

[Install]
WantedBy=xorg.target
