# user/xmonad.target

[Unit]
Description=Xmonad Window Manager
Requires=dbus.socket
Requires=xorg.target
AllowIsolate=true
