# user/xmonad.target

[Unit]
Description=Xmonad Window Manager
Requires=dbus.socket
AllowIsolate=true
