# user/xmonad.target
[Unit]
Description=Xmonad Window Manager
Requires=dbus.socket
Requires=xorg.target
BindsTo=xmonad@0.service
AllowIsolate=true
