# user/xmonad.service
[Unit]
Description=xmonad window manager

Conflicts=awesome.service openbox.service

After=xorg.service vncserver.service compton.service
After=xrandr-single.service xrandr-mirror.service xrandr-dual.service

Requisite=graphical.target
After=graphical.target

Requires=tint2.service
Before=tint2.service

Requires=wm.target
Before=wm.target

[Service]
Type=simple
SyslogIdentifier=xmonad
Slice=wm.slice

ExecStart=/usr/bin/xmonad

[Install]
WantedBy=desktop.target
WantedBy=laptop.target
WantedBy=tablet.target
WantedBy=vncserver.target
