# user/xmonad.target
[Unit]
Description=xmonad desktop environment
Conflicts=console.target

BindsTo=xmonad.service
After=xmonad.service
