# user/xorg.target
[Unit]
Description=X.org
Requires=dbus.socket
AllowIsolate=true

[Install]
WantedBy=awesome.target
WantedBy=openbox.target
