# user/xorg@.service
[Unit]
Description=Xorg server
StopWhenUnneeded=true

Wants=xserver@%i.target
Before=xserver@%i.target
Conflicts=x11vnc@%i.service vnc@%i.service xephyr@%i.service

[Service]
Type=notify
SyslogIdentifier=xorg
ExecStart=/usr/bin/xorg-launch-helper :%i -nolisten tcp -noreset vt7

