# user/xournal.service
[Unit]
Description=xournal

Requisite=wm.target
After=wm.target

[Service]
Type=simple
SyslogIdentifier=xournal
Slice=dropdown.slice

ExecStart=/usr/bin/xournal "${AB2_WORK_DIR}/notes/scratchpad.xoj"
