# user/xournal.service
[Unit]
Description=xournal

Requisite=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
Slice=apps.slice

ExecStart=/usr/bin/xournal "${AB2_MATHS_DIR}/notes/scratchpad.xoj"
