# user/xournal.service
[Unit]
Description=xournal

Requisite=graphical.target
After=graphical.target

[Service]
Type=simple
Slice=dropdown.slice

ExecStart=/usr/bin/xournal "${AB2_MATHS_DIR}/notes/scratchpad.xoj"
