# user/xrandr@.service
[Unit]
Description=X screen setup with randr
After=xorg.service
BindsTo=xorg.service
PropagatesReloadFrom=xorg.service

[Service]
Type=oneshot
ExecStart=%h/.config/systemd/scripts/screen.sh %i

[Install]
WantedBy=xorg.target
