# user/xrandr@.service
[Unit]
Description=X screen setup with randr
After=xorg@%i.service
BindsTo=xorg@%i.service

[Service]
Type=oneshot
Environment=DISPLAY=:%i
ExecStart=/home/abdo/.config/systemd/scripts/xrandr

[Install]
WantedBy=xorg.target
