# user/xrandr-single.service
[Unit]
Description=xrandr single
After=graphical.target
Requisite=graphical.target

[Service]
Type=oneshot
SyslogIdentifier=xrandr
ExecStart=/home/abdo/.config/systemd/scripts/xrandr.sh single
