# user/xrandr-xinerama@.service
[Unit]
Description=X screen setup: randr xinerama
After=xorg@%i.service
BindsTo=xorg@%i.service

[Service]
Type=oneshot
Environment=DISPLAY=:%i
ExecStart=/home/abdo/.config/systemd/scripts/xrandr xinerama

[Install]
WantedBy=xorg.target
