# user/xrandr@.service
[Unit]
Description=xrandr %i
After=xorg@0.service
Requisite=xorg@0.service

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/home/abdo/.config/systemd/scripts/xrandr %i
