# user/xrdb.service
[Unit]
Description=X resource setting
Documentation=man:xrdb(1)
After=xorg.service
BindsTo=xorg.service
PropagatesReloadFrom=xorg.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/systemd/conf/xenvironment
ExecStart=/usr/bin/xrdb -merge %h/.config/systemd/conf/Xresources
ExecStart=/usr/bin/xrdb -merge %h/.config/systemd/conf/Xresources.d/%H

[Install]
WantedBy=xorg.target
