# user/xrdb@.service
[Unit]
Description=X resource setting
Documentation=man:xrdb(1)
After=xorg@%I.service
BindsTo=xorg@%I.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/systemd/conf/xenvironment
Environment=DISPLAY=:%I
ExecStart=/usr/bin/xrdb -merge %h/.config/systemd/conf/Xresources
ExecStart=/usr/bin/xrdb -merge %h/.config/systemd/conf/Xresources.d/%H

[Install]
WantedBy=xorg.target
