# user/xrdb.service
[Unit]
Description=X resource setting
Documentation=man:xrdb(1)

PartOf=xserver.target
Before=xserver.target
After=xorg.service xwayland.service x11vnc.service vncserver.service

[Service]
Type=oneshot

ExecStart=/usr/bin/xrdb -merge %h/etc/xorg/Xresources
ExecStart=/usr/bin/xrdb -merge %h/etc/xorg/Xresources.d/%H

[Install]
WantedBy=graphical-session.target
