# user/xserver.target
[Unit]
Description=xwindow server

RefuseManualStart=true
StopWhenUnneeded=true

Conflicts=console.target
