# user/xserver@.target
[Unit]
Description=X server runing at that display
Requires=dbus.socket
StopWhenUnneeded=true
AllowIsolate=true
