# user/xserver@.target
[Unit]
Description=X server

StopWhenUnneeded=true
AllowIsolate=true
RefuseManualStart=true
