# user/xserver@.target
[Unit]
Description=X server
Requires=dbus.socket
StopWhenUnneeded=true
AllowIsolate=true
RefuseManualStart=true
