# shell/system/profile.d/profile.sh

# Signal whether profile is sourced
PROFILE_SOURCED=1
